# community/amazon-ssm-agent/amazon-ssm-agent.initd
#!/sbin/openrc-run
description="Initialize Amazon SSM Agent"
command="/usr/bin/amazon-ssm-agent"
command_background="true"
directory="/etc/amazon/ssm"
pidfile="/run/$RC_SVCNAME.pid"
retry="15"

depend() {
	need net
}
