# testing/nymphcast/nymphcast.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/bin/nymphcast_server"
command_args="-c /etc/nymphcast/nymphcast_config.ini -a /usr/share/nymphcast/apps/ -w /usr/share/nymphcast/wallpapers/"
