# testing/satellite/satellite.initd
#!/sbin/openrc-run
name="satellite"
description="satellite gemini service"
supervisor=supervise-daemon
command="/usr/bin/$name"
command_user="satellite:satellite"
pidfile="/run/$name.pid"
start_stop_daemon_args="--stdout $satellite_log_file --stderr $satellite_log_file"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f -m 0644 -o ${SVCNAME}:${SVCNAME} "$satellite_log_file"
}
