# testing/ganglia/gmetad.initd
#!/sbin/openrc-run

description="Ganglia monitoring daemon"
command="/usr/sbin/${RC_SVCNAME}"
command_args="--pid-file=/run/${RC_SVCNAME}.pid"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
	after firewall
}
