# testing/clash/clash.initd
#!/sbin/openrc-run

command="/usr/bin/clash"
command_args="-d /etc/clash"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
  need net
  use dns logger netmount
}
