# main/cyrus-sasl/saslauthd.initd
#!/sbin/openrc-run
command="/usr/sbin/saslauthd"
command_args=${SASLAUTHD_OPTS:--a sasldb}
pidfile="/run/saslauthd/saslauthd.pid"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory --mode 0775 /run/saslauthd
}
