# testing/ipp-usb/ipp-usb.initd
#!/sbin/openrc-run

name="ipp-usb"
description="Daemon for IPP over USB printer support"
command="/usr/sbin/ipp-usb"
command_args="standalone"
supervisor=supervise-daemon

depend() {
	need cupsd
	after avahi-daemon
}
