# testing/openswan/ipsec.initd
#!/sbin/runscript

# Openswan ipsec init.d file for alpine linux.

name=ipsec
daemon=/usr/libexec/ipsec/setup
pidfile=/var/run/pluto/ipsec_setup.pid

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting ${name}"
		$daemon start
	eend $?
}

stop() {
	ebegin "Stopping ${name}"
		$daemon stop
	eend $?
}

