# testing/irccd/irccd.initd
#!/sbin/openrc-run

name="IRC Client Daemon"
command="/usr/bin/irccd"
command_args="$IRCCD_OPTS"
pidfile="/var/run/irccd.pid"
start_stop_daemon_args="-b -m -u irccd -g irccd"
