# testing/triggerhappy/triggerhappy.initd
#!/sbin/runscript

: ${socket:=/var/run/thd.socket}
: ${triggers:=/etc/triggerhappy/triggers.d}

command=/usr/sbin/thd
pidfile=/var/run/thd.pid
command_args="--daemon --socket $socket --triggers $triggers $devices --pidfile $pidfile"

depend() {
	before udev
}
