# community/nftlb/nftlb.initd
#!/sbin/openrc-run

description="nftables load balancer"

command="/usr/sbin/nftlb"
command_args="
	${log_level:+--log $log_level}
	${cfgfile:+--config $cfgfile}
	${auth_key:+--key '$auth_key'}
	${listen_ipv6:+--ipv6}
	${listen_host:+--host $listen_host}
	${listen_port:+--port $listen_port}
	$command_args
	"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
