# testing/gpm/gpm.initd
#!/sbin/openrc-run

pidfile=/var/run/gpm.pid
command=/usr/sbin/gpm
command_args=${gpm_opts:--m /dev/input/mice -t imps2}

depend() {
	after localmount
	use hotplug logger
}
