# testing/xl2tpd/xl2tpd.initd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting xl2tpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/xl2tpd -- ${XL2TPD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping xl2tpd"
	start-stop-daemon --stop --quiet --pidfile /var/run/xl2tpd.pid
	result=$?
	start-stop-daemon --stop --quiet --exec /usr/sbin/xl2tpd-control
	result=$(( $result + $? ))
	eend $result
}
