# main/iperf3/iperf3.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

command="/usr/bin/iperf3"
command_args="--server $command_args"
description="A tool for performing network throughput measurements"

depend() {
	need net
	after firewall
}
