# main/unfs3/unfs3.initd
#!/sbin/openrc-run

depend() {
	config /etc/exports
	use net
	after firewall
	need portmap
}

start() {
	ebegin "Starting unfs3"
	start-stop-daemon --start --quiet --exec /usr/sbin/unfsd \
		-- ${UNFSD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping unfs3"
	start-stop-daemon --stop --quiet --exec /usr/sbin/unfsd
	eend $?
}

