# testing/pihole/pihole-gravity.timer
[Unit]
Description=Weekly ad-serving domains gathering

[Timer]
OnCalendar=weekly
AccuracySec=3h
Persistent=true
