# testing/pihole/pihole-logtruncate.timer
[Unit]
Description=Daily reset of dnsmasq/pi-hole query log

[Timer]
OnCalendar=daily
AccuracySec=5m
Persistent=true
