# community/tempo/tempo.initd
#!/sbin/openrc-run

configfile="${TEMPO_CONF:-/etc/tempo/tempo.yaml}"

command=/usr/bin/tempo
command_args="-config.file=${configfile}
	${TEMPO_OPTS}"
command_background=yes
command_user=tempo:grafana
make_pidfile=yes
pidfile=/var/run/tempo.pid
