# testing/ympd/ympd.initd
#!/sbin/runscript

name=ympd
command=/usr/bin/$name
command_args="-u $ympd_user -h $ympd_host -p $ympd_port -w $ympd_webport"
command_background=YES
pidfile=$ympd_pid

depend() {
        use logger dns
        need net
        after firewall
}
