# testing/prometheus-wireguard-exporter/prometheus-wireguard-exporter.initd
#!/sbin/openrc-run

description="Prometheus Wireguard Exporter"
supervisor=supervise-daemon
command=/usr/bin/prometheus_wireguard_exporter
command_args="-n /etc/wireguard/${INTERFACE}.conf"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
}
