# testing/lxd-lts/lxd-lts.initd
#!/sbin/openrc-run
# (c) 2016 Alpine Linux

command="/usr/sbin/lxd"
command_args="${LXD_OPTIONS}"
command_background="true"
pidfile="/var/run/${RC_SVCNAME}.pid"

depend() {
	need net dbus
	use lxcfs
	after firewall
}
