# community/dunst/dunst.initd
#!/sbin/openrc-run

description="dunst daemon"

command="/usr/bin/dunst"
supervisor=supervise-daemon
