# testing/signal-tlsd/signal-tlsd.initd
#!/sbin/openrc-run

name=signal-tlsd
description="signal-tlsd proxy daemon"
command=/usr/bin/signal-tlsd
command_args="$signal_tlsd_args"
command_user="$signal_tlsd_user:$signal_tlsd_group"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

error_logger="logger -t '${RC_SVCNAME}' -p daemon.info"

depend() {
	use dns logger
	need localmount net
	after firewall
}
