# testing/gradm/grsec-rbac.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Enabling grsecurity RBAC policy"
		grsec -E
	eend $?
}

stop() {
	ebegin "Disabling grsecurity RBAC policy"
		grsec -D
	eend $?
}

