# testing/mumble/murmur.initd
#!/sbin/runscript

depends() {
	needs localmount
}

description="Mumble VoIP Server"
command="/usr/bin/murmurd"
command_args="-ini /etc/murmur.ini"
pidfile="/var/run/murmur/murmur.pid"

# uses the openrc templates for start()/stop()

