# testing/razercfg/razerd.initd
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/razerd.pid"
command="/usr/bin/razerd"
command_args="-P $pidfile"
command_background="yes"

depend() {
	use logger
}
