# community/certmonger/certmonger.initd
#!/sbin/openrc-run
description="Certificate monitor and PKI enrollment client"
pidfile=/run/certmonger.pid
command="/usr/bin/certmonger"
command_args="-S -p $pidfile $CM_OPTS"

depend() {
	need dbus
	use net
}

