# community/swayidle/swayidle.initd
#!/sbin/openrc-run

command="/usr/bin/swayidle"
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
extra_started_commands="idle"

idle() {
	supervise-daemon $RC_SVCNAME --signal USR1
}
