# community/sane/saned.initd
#!/sbin/openrc-run

description="SANE network scanner server"

owner=saned
pidfile=/run/saned/saned.pid
cfgfile=/etc/sane.d/saned.conf
command=/usr/sbin/saned
command_args="-a $owner"
required_files="$cfgfile"

start_pre() {
	checkpath -dm755 -o $owner ${pidfile%/*}
}
