# testing/dpns/dpns.initd
#!/sbin/openrc-run

description="Docker Pull & Save"

depend() {
    need docker localmount net
}

start() {
    ebegin "Loading docker images"
    dpns-run
    eend $?
}
