# testing/pigpio/pigpiod.initd
#!/sbin/openrc-run

name="pigpiod"

command="/usr/bin/pigpiod"
command_args="-g -k"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"
