# main/net-snmp/snmptrapd.initd
#!/sbin/runscript

NAME=snmptrapd
DAEMON=/usr/sbin/$NAME

depend() {
	use logger
	need net
	after firewall
}

start() {
	ebegin "Starting ${NAME}"
		start-stop-daemon --start --quiet --background \
			--exec ${DAEMON} -- -p /var/run/${NAME}.pid ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"
		start-stop-daemon --stop --quiet --pidfile /var/run/${NAME}.pid
	eend $?
}
