# testing/kanidm/kanidm-unixd-tasks.initd
#!/sbin/openrc-run

name=kanidm-unixd-tasks
command="/usr/bin/kanidm_unixd_tasks"
command_user="root"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
