# testing/kanidm/kanidm-unixd.initd
#!/sbin/openrc-run

name=kanidm-unixd
command="/usr/bin/kanidm_unixd"
command_user="kanidm-unixd:kanidm-unixd"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
