# testing/kanidm/kanidmd.initd
#!/sbin/openrc-run

name=kanidmd
command="/usr/bin/kanidmd"
command_args="server -c /etc/kanidm/server.toml"
command_user="kanidmd:kanidmd"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
