# community/sanoid/sanoid.cron
#!/bin/sh
TZ=UTC sanoid --take-snapshots
