# community/sanoid/sanoid-prune.cron
#!/bin/sh
TZ=UTC sanoid --prune-snapshots
