# community/sanoid/sanoid.cron
#!/bin/sh
TZ=UTC /usr/bin/sanoid --cron
