# main/amavisd-new/amavisd.init
#!/sbin/runscript

PIDFILE="/var/amavis/amavisd.pid"

depend() {
	need net
	use logger
	use antivirus
	before mta
}

start() {
	ebegin "Starting amavisd-new"
	start-stop-daemon --start --quiet --name amavisd --pidfile ${PIDFILE} \
	 --exec /usr/sbin/amavisd
	eend $?
}

stop() {
	ebegin "Stopping amavisd-new"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
}
