# main/amavisd-new/amavisd.initd
#!/sbin/openrc-run

pidfile="/run/$RC_SVCNAME.pid"

name="amavisd"
command="/usr/sbin/$name"
command_args="-c /etc/amavisd.conf foreground"
command_background="yes"

start_stop_daemon_args="--user $user:$group"

depend() {
	need net
	use logger
	use antivirus
	before mta
}
