# community/networkmanager/networkmanager-dispatcher.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/libexec/nm-dispatcher
command_args_foreground="--persist"

description="Network Manager Dispatcher Daemon"

depend() {
   need dbus
   before networkmanager
}
