# testing/olsrd/olsrd.initd
#!/sbin/openrc-run

name=olsrd
command="/usr/sbin/$name"
command_args="$olsrd_opts"
command_background="yes"

start_stop_daemon_args="-g wheel -k 027"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
