# testing/repowerd/repowerd.initd
#!/sbin/openrc-run

name="Repowerd"
description="Power daemon to monitor and control system power state"

command=/usr/sbin/repowerd
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need udev
}
