# community/wireplumber/wireplumber.initd
#!/sbin/openrc-run

description="Multimedia processing graphs"

command="/usr/bin/wireplumber"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depends() {
	use dbus
}
