# community/iwd/ead.initd
#!/sbin/openrc-run

description="Ethernet authentication daemon"

: ${error_logger="logger -t ead -p daemon.info"}

command="/usr/libexec/ead"
command_background="yes"
command_args="
	${allowed_interfaces:+"--interfaces \"$allowed_interfaces\""}
	${ignored_interfaces:+"--nointerfaces \"$ignored_interfaces\""}
	${command_args:-}
	"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need dbus
	before net
	keyword -shutdown
}
