# community/iwd/iwd.initd
#!/sbin/openrc-run

description="iNet wireless daemon"

: ${error_logger="logger -t iwd -p daemon.info >/dev/null 2>&1"}
: ${start_wait=50}  # milliseconds

command="/usr/libexec/iwd"
command_background="yes"
command_args="
	${allowed_interfaces:+"--interfaces \"$allowed_interfaces\""}
	${ignored_interfaces:+"--nointerfaces \"$ignored_interfaces\""}
	${command_args:-}
	"
start_stop_daemon_args="
	${start_wait:+--wait $start_wait}
	${start_stop_daemon_args:-}
	"
supervise_daemon_args=""
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need dbus
	before net
	keyword -shutdown
	provide wlan
}
