# testing/tqftpserv/tqftpserv.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name="tqftpserv"
description="Qualcomm Trivial File Transfer Protocol Server"
command="/usr/bin/tqftpserv"

depend() {
	before rmtfs
	use qrtr-ns
}
