# community/sndio/sndiod.initd
#!/sbin/openrc-run

# Sample init.d file for alpine linux.

name="Sndio server"
pname="sndiod"
user="sndiod"
group="sndiod"
command="/usr/sbin/$pname"
command_args="$DAEMON_OPTS"
command_background="yes"

start_stop_daemon_args="--user $user:$group"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
