# testing/greetd/greetd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Greeter daemon"
description="Greeter daemon"

command=/usr/bin/greetd

depend() {
    provide display-manager
    need localmount
}

start_pre() {
    checkpath --owner greetd:greetd --directory /var/empty
}
