# community/fnott/fnott.initd
#!/sbin/openrc-run

command="/usr/bin/fnott"
supervisor=supervise-daemon
command_args="--print-pid=3 $command_args"
notify=fd:3

depend() {
	need dbus
}
