# main/gradm/grsec-rbac.initd
#!/sbin/runscript

start() {
	ebegin "Enabling grsecurity RBAC policy"
		gradm -E
	eend $?
}

stop() {
	ebegin "Disabling grsecurity RBAC policy"
		gradm -D
	eend $?
}

