# testing/gearmand/gearmand.initd
#!/sbin/openrc-run

name="gearmand"
description="Gearman server"

command="/usr/sbin/gearmand"
command_args="-d $gearmand_opts"
command_background="yes"

pidfile="/run/${RC_SVCNAME}.pid"
start_stop_daemon_args="--user=$gearmand_user:$gearmand_group"

depend() {
	need net localmount
	after firewall
}
