# testing/remco/remco.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

command="/usr/bin/remco"
command_background="yes"
error_log="/var/log/remco.log"

pidfile="/var/run/${SVCNAME}.pid"

depend() {
	need net
	after firewall
}
