# community/ympd/ympd.initd
#!/sbin/openrc-run

name=ympd
command=/usr/bin/$name
command_args="-u $ympd_user"
command_background=YES
pidfile=$ympd_pid

depend() {
        use logger dns
        need net
        after firewall
}

start_pre() {
        checkpath --directory --owner "$ympd_user" \
		--mode 0775 /var/run/"$name"
}
