# unmaintained/nldev/nldev-trigger.initd
#!/sbin/openrc-run

depend() {
	need nldev
	provide dev
	keyword -lxc -vserver
}

start() {
	ebegin "Populating /dev with existing devices through uevents"
	/bin/nltrigger add all
	eend $?
}
