# unmaintained/lizardfs/lizardfs-cgiserv.initd
#!/sbin/openrc-run

pidfile=/var/run/lizardfs-cgiserv.pid
command=/usr/sbin/lizardfs-cgiserver


depend() {
	after net
}

