# testing/minisatip/minisatip.initd
#!/sbin/openrc-run
command="/usr/bin/minisatip"
command_args="-f $extra_opts"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
command_user="minisatip:minisatip"
capabilities="^cap_net_bind_service"

depend() {
	need net
	use logger dns
	after firewall
}
