# community/stubby/stubby.initd
#!/sbin/openrc-run

name="Stubby DNS Privacy stub resolver"

command="/usr/bin/stubby"
command_args="$command_args -C /etc/stubby/stubby.yml"
command_background=true
pidfile="/run/${RC_SVCNAME}"
command_user="stubby:stubby"
capabilities="^cap_net_bind_service"

depend() {
	need net
	after firewall
}
