# main/icecast/icecast.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting Icecast 2"
	start-stop-daemon --start --exec /usr/bin/icecast \
		-- ${ICEC_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping Icecast 2"
	start-stop-daemon --stop --exec /usr/bin/icecast
	eend $?
}
