# testing/zerotier-one/zerotier-one.initd
#!/sbin/openrc-run

depend() {
    need net
}

command="/usr/sbin/zerotier-one"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
