# testing/zerotier-one/zerotier-one.initd
#!/sbin/openrc-run

command="/usr/sbin/zerotier-one"
command_args="$ZEROTIER_OPTS"
command_background="yes"

start_stop_daemon_args=""
pidfile="/var/lib/zerotier-one/zerotier-one.pid"

depend() {
	after net
	before firewall netmount
}
