# testing/zerotier-one/zerotier-one.initd
#!/sbin/openrc-run

depend() {
    provide net
    use network
    before firewall netmount
}

command="/usr/sbin/zerotier-one"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
