# main/dbus/dbus-user.initd
#!/sbin/openrc-run

name="Login Message Bus"
description="D-Bus Login Session Message Bus"

supervisor=supervise-daemon
command="/usr/bin/dbus-daemon"
command_args="--session --nofork --nopidfile --address=$DBUS_SESSION_BUS_ADDRESS"
