# testing/phabricator/phd.initd
#!/sbin/openrc-run

command="/usr/share/php/phabricator/bin/phd"
command_args="start"
start_stop_daemon_args="--quiet"

stop() {
        $command stop
}

