# community/mce/mce.initd
#!/sbin/openrc-run

name=mce
description="Mode Control Entity (MCE)"
command="/usr/sbin/mce"
command_background=yes

depend() {
	need dbus
}
