# testing/frr/frr
#!/sbin/openrc-run

description="FRRouting"

start() {
        exec /usr/lib/frr/frr start
}

stop() {
        exec /usr/lib/frr/frr stop
}

restart() {
        exec /usr/lib/frr/frr restart
}

status() {
        exec /usr/lib/frr/frr status
}
