# testing/frr/frr.initd
#!/sbin/openrc-run

description="FRRouting"
extra_started_commands="reload"

start() {
        /usr/lib/frr/frrinit.sh start
}

stop() {
        /usr/lib/frr/frrinit.sh stop
}

restart() {
        /usr/lib/frr/frrinit.sh restart
}

status() {
        /usr/lib/frr/frrinit.sh status
}

reload() {
	/usr/lib/frr/frrinit.sh reload
}
