# community/milter-greylist/milter-greylist.initd
#!/sbin/openrc-run

command=/usr/bin/$SVCNAME
pidfile=/var/run/$SVCNAME/$SVCNAME.pid

depend() {
	need net
	after firewall
}
