# community/ntfy-alertmanager/ntfy-alertmanager.initd
#!/sbin/openrc-run

name=ntfy-alertmanager
command="/usr/bin/ntfy-alertmanager"
command_user="ntfy-alertmanager:ntfy-alertmanager"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
