# main/collectd/collectd.initd
#!/sbin/openrc-run

DAEMON=/usr/sbin/collectd

depend() {
	need net
}

start() {
	ebegin "Starting collectd"
	start-stop-daemon --start --exec $DAEMON -- $collectd_options
	eend $?
}

stop () {
	ebegin "Stopping collectd"
	start-stop-daemon --stop --exec $DAEMON
	eend $?
}

