# community/guacamole-server/guacd.initd
#!/sbin/openrc-run

name="Guacamole Server"

: ${command_user:="guacd"}

command="/usr/sbin/guacd"
command_args="-f $command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
