# testing/grommunio-gromox/gromox-event.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
description="Gromox event service"
command="/usr/libexec/gromox/event"

depend() {
    need mariadb
}

start_pre() {
    checkpath -d -m 750 -o gromox:gromox /run/gromox
}
