# testing/grommunio-gromox/gromox-pop3.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
description="Gromox POP3 server"
description_reload="Reload configuration without exiting"

command="/usr/libexec/gromox/pop3"

extra_started_commands="reload"

depend() {
    after mariadb
}

reload() {
    ebegin "Reloading configuration"
    $supervisor $RC_SVCNAME --signal HUP
    eend $?
}
