# testing/grommunio-gromox/gromox-zcore.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/libexec/gromox/zcore"

depend() {
    need mariadb atd
}

start_pre() {
    checkpath -d -m 750 -o gromox:gromox /run/gromox
}
