# testing/grommunio-gromox/gromox-pop3.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/libexec/gromox/pop3"

depend() {
    need mariadb
}

start_pre() {
    checkpath -d -m 750 -o gromox:gromox /run/gromox
}
