# testing/rustdesk-server/hbbr.initd
#!/sbin/openrc-run

name=hbbr
description="Rust-desk server hbbr daemon"
command="/usr/bin/hbbr"
command_background=true
command_args=${command_args:-"-k _"}
: ${directory:=/var/lib/rustdesk-server}
: ${command_user:=rustdesk:rustdesk}
: ${umask:=0022}

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -o $command_user "$directory"
}
