# testing/rustdesk-server/hbbs.initd
#!/sbin/openrc-run

name=hbbs
description="Rust-desk server hbbs daemon"
command="/usr/bin/hbbs"
command_background=true
: ${directory:=/var/lib/rustdesk-server}
: ${command_user:=rustdesk:rustdesk}
: ${umask:=0022}

depend() {
	need net hbbr
	after firewall
}

start_pre() {
	checkpath -d -o $command_user "$directory"
}
