# community/csync2/csync2.initd
#!/sbin/openrc-run

command="/usr/sbin/csync2"
command_args="-ii ${command_args:-}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
