# testing/hvtools/hv_kvp_daemon.initd
#!/sbin/openrc-run

name="Hyper-V daemon: hv_kvp_daemon"
command=/usr/bin/hv_kvp_daemon

depend() {
	use clock logger
	need localmount net
}

start_pre() {
	# Delete the existing store
	rm -rf /var/lib/hyperv
}
	
