# testing/runit/runit.initd
#!/sbin/openrc-run

command="/sbin/runsvdir"
command_args="/etc/service"
pidfile="/run/${SVCNAME}.pid"
command_background="true"
description="starts and monitors a collection of runsv(8) processes"
retry="SIGHUP/5"

depend() {
    need localmount
    after firewall
}

