# testing/lldpd/lldpd.initd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name=lldpd
pidfile=/run/lldpd/lldpd.pid
command=/usr/sbin/lldpd
retry="TERM/10/KILL/5"
command_args="${LLDPD_OPTS}"
start_stop_daemon_args="--wait 300"

depend() {
	use net
}

start_pre() {
	checkpath --directory /run/lldpd -o lldpd -m 0700
}

