# testing/salt/salt-syndic.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/salt/files/syndic-initd-1,v 1.1 2012/04/24 23:37:27 chutzpah Exp $

command="/usr/bin/salt-syndic"
command_args="--daemon ${SALT_OPTS}"
pidfile="/var/run/salt-syndic.pid"
name="SALT syndic daemon"

depend() {
	use net
}
