# testing/salt/salt-minion.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/salt/files/minion-initd-1,v 1.1 2012/04/24 23:37:27 chutzpah Exp $

command="/usr/bin/salt-minion"
command_args="--daemon ${SALT_OPTS}"
pidfile="/var/run/salt-minion.pid"
name="SALT minion daemon"

depend() {
	need net
}
