# testing/fileshelter/fileshelter.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="fileshelter"
command="/usr/bin/fileshelter"
command_user="fileshelter:fileshelter"
command_background=true
pidfile="/run/fileshelter.pid"

depend() {
	need net localmount
	after firewall
}
