# testing/pimd/pimd.initd
#!/sbin/openrc-run

# PIMD init.d file for alpine linux.

name=pimd
command="/usr/sbin/$name"
command_args="$pimd_opts"
command_background="yes"

pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
