# community/corerad/corerad.initd
#!/sbin/openrc-run

name=CoreRAD
description="Extensible and observable IPv6 Neighbor Discovery Protocol router advertisement daemon"
command=/usr/sbin/corerad
command_args="-c /etc/corerad/config.toml"
command_user="corerad:corerad"
supervisor="supervise-daemon"
capabilities="^cap_net_admin,^cap_net_raw"

depend() {
        need net
        after firewall
}
