# community/keyd/keyd.initd
#!/sbin/openrc-run

name="keyd"
description="A key remapping deamon"

: ${error_logger="logger -t $name -p daemon.warning >/dev/null 2>&1"}

command="/usr/bin/keyd"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need localmount dev-mount
}
