# testing/grommunio-index/grommunio-index.cron
#!/bin/sh
# Recurrent regeneration of grommunio FTS indexes

/usr/bin/grommunio-index -Aq > /dev/null 2>&1
