# community/x2goserver/x2goserver.initd
#!/sbin/openrc-run

DAEMON=/usr/sbin/x2gocleansessions
PIDFILE=/var/run/x2goserver.pid
 
depend() {
	need net
	need localmount
 	after firewall
}
 
start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --quiet --pidfile ${PIDFILE} --exec ${DAEMON} -- $DAEMON-OPTS
	eend $?
}
 
stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
}
 
