# community/modemmanager/modemmanager.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/sbin/ModemManager
command_args="${modemmanager_opts}"

description="ModemManager Daemon"

depend() {
   need dbus
}
