# testing/haproxy-dataplaneapi/haproxy-dataplaneapi.initd
#!/sbin/openrc-run

name="haproxy-dataplaneapi"
description="HAProxy Data Plane API"

pidfile="/run/$RC_SVCNAME.pid"
command="/usr/sbin/dataplaneapi"
command_args="-f $CONFIGFILE"
command_background="yes"
command_user="haproxy:haproxy"

depend() {
	need net
	after firewall
}
