# main/lvm2/lvmlockd.initd
#!/sbin/openrc-run

name=lvmlockd
command="/sbin/$name"
command_args_foreground="--foreground"
command_background=true
pidfile="/run/$name.pid"

start_pre() {
    checkpath --directory /run/lvm
}
