# testing/nullmailer/nullmailer.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=nullmailer
description="Nullmailer relay-only MTA"

command=/usr/sbin/nullmailer-send

depend() {
	need net
}
