# testing/virtualbox/vboxweb.initd
#!/sbin/openrc-run

description="VirtualBox web service"

pidfile="/var/run/vboxweb-service.pid"
command="/usr/bin/vboxwebsrv"
command_args="--pidfile ${pidfile}"
start_stop_daemon_args="--make-pidfile --pidfile ${pidfile} --background"

depend() {
        need networking
}
