# community/sddm/sddm.initd
#!/sbin/openrc-run

name="Display Manager"
description="SDDM Display Manager"

command=/usr/bin/sddm
command_args="$sddm_opts"

depend() {
	provide display-manager
	need localmount
}

start_pre() {
	checkpath -d -o sddm:video -m755 /var/lib/sddm
}
