# testing/rng-tools/rngd.initd
#!/sbin/openrc-run

name=rngd
daemon=/usr/sbin/$name

depend() {
	need localmount
	after urandom
	provide entropy
}

start() {
	ebegin "Starting ${name}"
	start-stop-daemon --start \
		--pidfile $pidfile --exec $command -- $command_args
        eend $?
}

stop() {
	ebegin "Stopping ${name}"
	if [ -f $pidfile ]; then
		start-stop-daemon --stop --pidfile $pidfile
	fi
	ewend $?
}

