# community/prometheus-nextcloud-exporter/nextcloud-exporter.initd
#!/sbin/openrc-run
name="prometheus-nextcloud-exporter"
supervisor=supervise-daemon
command="/usr/bin/nextcloud-exporter"
command_args="$ARGS"
command_background="yes"
command_user="prometheus:prometheus"
pidfile="/run/${RC_SVCNAME}.pid"
