# community/iio-sensor-proxy/iio-sensor-proxy.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/libexec/iio-sensor-proxy"

depend() {
	need dbus
}
