# community/sway/sway.initd
#!/sbin/openrc-run

command="/usr/bin/sway"
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

# FIXME export WAYLAND_DISPLAY SWAYSOCK XDG_CURRENT_DESKTOP=sway
