# community/gdm/gdm.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Display Manager"
description="GNOME Display Manager"

command=/usr/sbin/gdm
command_args="$GDM_OPTS"

depend() {
	provide display-manager
	need net dbus
}
