# community/nextcloud28/nextcloud28.cron
#!/bin/sh

# Run only when nextcloud service is started.
if rc-service nextcloud -q status >/dev/null 2>&1; then
	su nextcloud -s /bin/sh -c 'php82 -f /usr/share/webapps/nextcloud/cron.php'
fi
