# unmaintained/fancontroled/fancontroled.initd
#!/sbin/runscript

command=/usr/sbin/fancontroled
command_args="${FANCONTROLED_OPTS}"
pidfile=/var/run/fancontroled.pid
description="Temperature dependent fan speed control"

depend() {
	need modules
}

start_pre() {
	if [ -z "$FANCONTROLED_OPTS" ]; then
		eerror "Please set FANCONTROLED_OPTS in /etc/conf.d/$RC_SVCNAME"
		return 1
	fi
}

