# testing/netdata/netdata.initd
#!/sbin/openrc-run

command=/usr/sbin/netdata
command_args=""
pidfile=/run/netdata.pid
name="netdata"
description="Real-time performance monitoring"
