# community/frr/frr.initd
#!/sbin/openrc-run

description="FRRouting"

start() {
        /usr/lib/frr/frr start
}

stop() {
        /usr/lib/frr/frr stop
}

restart() {
        /usr/lib/frr/frr restart
}

status() {
        /usr/lib/frr/frr status
}
