# testing/coredns/coredns.initd
#!/sbin/openrc-run

name=coredns
description="CoreDNS is a fast and flexible DNS server"
command=/usr/bin/coredns
command_args="$ARGS"
command_user="coredns:coredns"
supervisor="supervise-daemon"

depend() {
        need net
	provide dns
}
