# testing/ntpd-rs/ntpd-rs.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Network Time Service (ntpd-rs)"
supervisor="supervise-daemon"
command="/usr/bin/ntp-daemon"
supervise_daemon_args=" -e "RUST_LOG=\"info"\""
command_user="ntpd-rs:ntpd-rs"
capabilities="^cap_sys_time"

depend() {
	need net
	provide ntp-client
	after firewall
}
