# testing/ntpd-rs/ntpd-rs-metrics-exporter.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Network Time Service (ntpd-rs) metrics exporter"
supervisor="supervise-daemon"
command="/usr/bin/ntp-metrics-exporter"
supervise_daemon_args=" -e "RUST_LOG=\"info"\""
command_user="ntpd-rs-observe:ntpd-rs-observe"

depend() {
	need net
	after firewall
}
