# community/captagent/captagent.initd
#!/sbin/openrc-run

# init.d file for captagent
# Copyright (c) 2017 Kaarle Ritvanen

name=captagent
command=/usr/sbin/captagent
command_args="-d"
pidfile=/var/run/$RC_SVCNAME.pid

depend() {
	need net
}
