# community/snapcast/snapcast-client.initd
#!/sbin/openrc-run

name=snapclient
command="/usr/bin/snapclient"
command_args="-d $snapclient_opts"
command_user="snapcast"
pidfile="/run/snapclient/pid"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory --owner snapcast:audio --mode 0775 \
		/run/snapclient
}

