# community/snapcast/snapcast-client.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name=snapclient
command="/usr/bin/snapclient"
command_args="$snapclient_opts"
command_user="snapcast"

depend() {
	need net
	after firewall
}
