# community/snapcast/snapcast-server.initd
#!/sbin/openrc-run

name=snapserver
command="/usr/bin/snapserver"
command_args="-d $snapserver_opts"
command_user="snapcast"
start_stop_daemon_args="--quiet"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory --owner snapcast:audio --mode 0775 \
		/run/snapserver
}

