# community/openvas-scanner/greenbone-nvt-sync.cron
#!/bin/sh
/usr/bin/greenbone-nvt-sync >& /dev/null
/usr/bin/greenbone-scapdata-sync >& /dev/null
/usr/bin/greenbone-certdata-sync >& /dev/null

