# community/openvas-scanner/greenbone-nvt-sync.cron
0  2 *  *       *       /usr/bin/greenbone-nvt-sync >& /dev/null
