# testing/sslh/sslh.initd
#!/sbin/openrc-run

: ${cfgfile:="/etc/sslh.conf"}
: ${mode:="fork"}
: ${wait:=50}  # milliseconds

description="Port multiplexer for SSH, HTTPS, OpenVPN etc."

command="/usr/sbin/sslh-$mode"
command_args="-F$cfgfile -f ${command_args:-$DAEMON_OPTS}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="--wait $wait"

required_files="$cfgfile"
