# testing/bitwarden_rs/bitwarden_rs.initd
#!/sbin/openrc-run

name="bitwarden_rs"
description="bitwarden_rs alternative bitwarden backend"
supervisor="supervise-daemon"
command="/usr/sbin/bitwarden_rs"
command_user="bitwarden_rs:bitwarden_rs"
output_log="/var/log/bitwarden_rs/access.log"
error_log="/var/log/bitwarden_rs/error.log"
directory=$DATA_FOLDER
export DATA_FOLDER

depend() {
	need localmount net
	after firewall
}
