# testing/lxdm/lxdm.initd
#!/sbin/runscript

depends() {
	needs localmount
}

description="Lightweight Login Manager"
command="/usr/sbin/lxdm"
command_args="-d"
pidfile="/var/run/lxdm.pid"

# uses the openrc templates for start()/stop()
