# testing/lxdm/lxdm.initd
#!/sbin/runscript

depends() {
	needs localmount
}

start() {
	ebegin "Starting lxdm"
	start-stop-daemon --start --quiet --exec /usr/sbin/lxdm -- -d
	eend $?
}

stop() {
	ebegin "Stopping lxdm"
	start-stop-daemon --stop --quiet --exec /usr/sbin/lxdm
	eend $?
}

