# community/syncthing13/syncthing.initd
#!/sbin/openrc-run

name=syncthing
command=/usr/bin/$name
command_args=$syncthing_args
start_stop_daemon_args="--background --user $name --group $name"

depend() {
	use logger dns
	need net
	after firewall
}

