# community/swaylock/swaylock.initd
#!/sbin/openrc-run

command="/usr/bin/swaylock"
supervisor=supervise-daemon
command_args="--ready-fd 3 $command_args"
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
notify=fd:3
extra_started_commands="unlock"

unlock() {
	supervise-daemon $RC_SVCNAME --signal USR1
}
