# testing/sopel/sopel.initd
#!/sbin/openrc-run

command="/usr/bin/sopel"
command_args="-c $CONFIG_PATH"
pidfile="/run/sopel.pid"
command_background=true
command_user="sopel:sopel"

depend() {
	need net
}
