# community/himitsu/himitsud.initd
#!/sbin/openrc-run

description="himitsu key ring daemon"

command="/usr/bin/himitsud"
command_args="-R 3"
supervisor=supervise-daemon
notify=fd:3
