# community/ofono/ofono.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="oFono"
description="oFono mobile telephony daemon"

command="/usr/sbin/ofonod"
command_args_foreground="-n"

depend() {
	need dbus
}
