# community/repmgr/repmgrd.initd
#!/sbin/openrc-run

command="/usr/bin/repmgrd"
command_background="yes"
# "user" variable is here for backward compatibility only
command_user="${command_user:-${user:-postgres}}"
pidfile="/run/$RC_SVCNAME.pid"
command_args="--no-pid-file --daemonize=false"

depend() {
	need postgresql
	use net logger
}
