# community/repmgr/repmgrd.initd
#!/sbin/openrc-run

command="/usr/bin/repmgrd"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="--user ${user:-postgres}"

depend() {
	need postgresql
	use net logger
}
