# testing/semaphoreui/semaphoreui.initd
#!/sbin/openrc-run

name="Semaphore UI"

: "${command_user:="semaphore:semaphore"}"

command="/usr/bin/semaphore"
: "${config_file:="/var/lib/semaphore/config.json"}"
: "${command_args:="server --config $config_file"}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

reload() {
	${supervisor} ${RC_SVCNAME} --signal HUP --pidfile "${pidfile}"
}
