# main/devicemaster-linux/nslink.initd
#!/sbin/openrc-run

program=/sbin/nslink
pidfile=/var/run/nslink.pid

start_pre() {
	if ! [ -f /etc/nslink.conf ]; then
		eerror "Configuration file /etc/nslink.conf not found"
		return 1
	fi

	if ! grep -q -w nslink /proc/modules; then
		modprobe nslink
	fi
}

