# testing/tuptime/tuptime.initd
#!/sbin/openrc-run

name="tuptime"
command="/usr/bin/tuptime"
command_user="_tuptime"
command_args="-q"

depend() {
	need localmount
	after bootmisc
	after clock
}

stop() {
	ebegin "Stopping ${name}"
	su ${command_user} -s /bin/sh -c "${command} ${command_args} -g"
	eend $?
}

status() {
	ebegin "Status ${name}"
	su ${command_user} -s /bin/sh -c "${command} ${command_args}"
	eend $?
}
