# testing/tuptime/tuptime.cron
#!/bin/sh
if [ -x /usr/bin/tuptime ]; then
	su tuptime -s /bin/sh -c '/usr/bin/tuptime -x' > /dev/null
fi
