# testing/tuptime/tuptime.cron
#!/bin/sh
if [ -x /usr/bin/tuptime ]; then /usr/bin/tuptime -x > /dev/null; fi
