# main/oscam/oscam.initd
#!/sbin/runscript

name=oscam
command=/usr/bin/oscam
command_args=$oscam_args
pidfile=$oscam_pid

depend() {
	use logger dns
	need net
	after firewall
}
