# main/oscam/oscam.initd
#!/sbin/openrc-run

name=oscam
command=/usr/bin/oscam
command_args=$oscam_args
pidfile=$oscam_pid
start_stop_daemon_args="${oscam_user:+--user} $oscam_user"

depend() {
	use logger dns
	need net
	after firewall
}

start_pre() {
	checkpath --directory --owner $oscam_user:video --mode 0775 \
		/var/run/$name
}
