# community/cacti/cacti.crontab
*/5 * * * * php82 /usr/share/webapps/cacti/poller.php >/dev/null 2>&1
