# community/cacti/cacti.crontab
*/5 * * * * php8 /usr/share/webapps/cacti/poller.php >/dev/null 2>&1
