# community/jool-tools/joold.initd
#!/sbin/openrc-run

: ${cfgfile:="/etc/jool/$RC_SVCNAME.conf"}

command="/usr/bin/joold"
command_args="$cfgfile"
command_background=yes

pidfile="/run/$RC_SVCNAME.pid"
required_files="$cfgfile"

depend() {
	need net jool
	after keepalived
}
