# testing/jitterentropy-rngd/jitterentropy-rngd.initd
#!/sbin/openrc-run

description="CPU Jitter Random number generator daemon"

pidfile=/run/$RC_SVCNAME.pid
command=/usr/sbin/jitterentropy-rngd
command_args="--pid $pidfile"

depend() {
	need localmount
	after urandom
	provide entropy
}
