# testing/crowdsec/crowdsec.initd
#!/sbin/openrc-run

command=/usr/bin/crowdsec
command_args="-c /etc/crowdsec/config.yaml"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
extra_started_commands="reload"
description_reload="Reload configuration"

depend() {
	after firewall
}

reload() {
	ebegin "Reloading $RC_SVCNAME"
	start-stop-daemon --signal HUP --pidfile "$pidfile"
	eend $?
}
