# testing/librespot/librespot.initd
#!/sbin/openrc-run

command=/usr/bin/librespot
command_args="${librespot_args}"
pidfile=/run/librespot.pid
command_background=true
logfile="/var/log/librespot.log"
start_stop_daemon_args="--user ${librespot_user:-librespot} --group ${librespot_group:-audio} --stdout $logfile --stderr $logfile"

depends() {
	need localmount net
	before firewall
}
