# testing/cockroach/cockroach.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=CockroachDB
description="The open source, cloud-native SQL database."

command=/usr/sbin/cockroach
command_args="${cockroach_opts}"
command_user="cockroach:cockroach"
pidfile=/run/${RC_SVCNAME}.pid
