# community/virtualbox-guest-additions/virtualbox-drm-client.initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="VBoxDRMClient control service"

command="/usr/sbin/VBoxDRMClient"
output_log="/var/log/VBoxDRMClient.log"
error_log="/var/log/VBoxDRMClient.log"
pidfile="/var/run/virtualbox-drm-client.pid"
command_background="yes"
