# community/docker/docker-user.initd
#!/sbin/openrc-run

name="Rootless Docker Daemon"
description="Persistent process that manages docker containers"

supervisor=supervise-daemon
command="dockerd-rootless"
command_args=""
