# testing/matridge/matridge.initd
#!/sbin/openrc-run

description="Matrix/XMPP gateway"
command=/usr/bin/matridge
: ${command_user:=matridge:matridge}
: ${command_args:="-c /etc/matridge/matridge.conf"}

output_logger="logger -t matridge -p daemon.info"
error_logger="logger -t matridge -p daemon.err"

supervisor="supervise-daemon"

depend() {
	need net
	use dns logger
	after firewall
}
