# testing/tvheadend/tvheadend.initd
#!/sbin/runscript
depend() {
        need net
}

start() {
        ebegin "Starting tvheadend" 
        start-stop-daemon --start \
                --user ${TVHUSER} \
                --exec ${TVHBIN} -- -C -f
        eend $?
}

stop() {
        ebegin "Stopping tvheadend" 
        start-stop-daemon --stop --quiet --name tvheadend --signal 2
        eend $?
}
