# testing/autofs/autofs.initd
#!/sbin/openrc-run

name=${RC_SVCNAME}

pidfile="/run/${RC_SVCNAME}.pid"

command="/usr/sbin/automount"
command_args="$automount_opts --pid-file $pidfile"
command_background="yes"

depend() {
  need net
  after firewall
}
