# testing/autofs/autofs.initd
#!/sbin/openrc-run

name=${RC_SVCNAME}

command="/usr/sbin/automount"
command_args="$automount_opts"
command_background="yes"

pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  after firewall
}
