# testing/libc6/nscd.initd
#!/sbin/runscript

DAEMON=/usr/sbin/nscd

start() {
        ebegin "Starting nscd"
        # remove stale files
        rm -f /var/db/nscd/* /var/run/nscd/* 2>/dev/null
        start-stop-daemon --start --exec $DAEMON
        eend $?
}

stop () {
        ebegin "Stopping nscd"
        start-stop-daemon --stop --quiet \
		--pidfile=/var/run/nscd/nscd.pid 
        eend $?
}
