# unmaintained/netdata/netdata.initd
#!/sbin/openrc-run

command=/usr/sbin/netdata
pidfile=/run/netdata.pid
command_args="-P ${pidfile}"
name="netdata"
description="Real-time performance monitoring"

depend() {
	need net
	after apache2 squid nginx mysql named opensips hostapd postfix lm_sensors
}
