# testing/freeswitch/freeswitch.initd
#!/sbin/runscript

NAME=freeswitch
USER=freeswitch
DAEMON=/usr/bin/freeswitch
OPTIONS_START="-u $USER -g $GROUP -nc"
OPTIONS_STOP="-stop"

start() {
	ebegin "Starting $NAME"
		$DAEMON $OPTIONS_START 
	eend 0
	}

stop() {
	ebegin "Stopping $NAME"
		$DAEMON $OPTIONS_STOP
	eend 0
	}
