# testing/otf/otf-agent.initd
#!/sbin/openrc-run

name="otf agent"

: ${command_user:="otf:otf"}

command="/usr/bin/otf-agent"
command_args="--token '$(cat $token_file)' $command_args"
command_background="yes"
error_log="/var/log/otf-agent.log"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

start_pre() {
        checkpath -f --owner "$command_user" --mode 0644 \
                $error_log
}
