# testing/otf/otfd.initd
#!/sbin/openrc-run

name="otf server"

: ${command_user:="otf:otf"}

command="/usr/bin/otfd"
command_args="--secret '$(cat $secret_file)' $command_args"
command_background="yes"
error_log="/var/log/otfd.log"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

start_pre() {
        checkpath -f --owner "$command_user" --mode 0644 \
                $error_log
}
