# testing/buildkit/buildkitd.initd
#!/sbin/openrc-run

name="buildkitd"
description="buildkitd builder daemon"
command="/usr/bin/buildkitd"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"
directory="/var/lib/buildkit"
error_log="/var/log/buildkitd.log"

depend() {
	need net localmount
	after firewall
}

start_pre() {
	checkpath -d "$directory"
	checkpath -f "$error_log"
}
