# testing/scsi-tgt/tgtd.initd
#!/sbin/openrc-run

name="tgtd"
command="/usr/sbin/tgtd"
command_args_foreground="--foreground"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"
# sigterm doesn't stop it
retry="SIGINT/5"

output_log="/var/log/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME.log"

depend() {
	need dev localmount
}
