# testing/kubernetes/kube-controller-manager.initd
#!/sbin/openrc-run
# Copyright 2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

supervisor=supervise-daemon
description="Kubernetes Controller Manager service"

if [ -e /var/lib/kubernetes/kube-controller-manager-flags.env ]; then
	. /var/lib/kubernetes/kube-controller-manager-flags.env;
fi

command="/usr/bin/kube-controller-manager"
command_args="${command_args} ${KUBE_CONTROLLER_MANAGER_ARGS}"
: ${output_log:=/var/log/$RC_SVCNAME/$RC_SVCNAME.log}
: ${error_log:=/var/log/$RC_SVCNAME/$RC_SVCNAME.log}

depend() {
	after net
}
