# main/vlc/vlc.initd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting VLC"
	start-stop-daemon --start --oknodo \
		--pidfile /var/run/vlc.pid --user vlc \
		--exec /usr/bin/vlc -- ${VLC_OPTS}
	eend $?
}

stop () {
	ebegin "Stopping VLC"
	start-stop-daemon --stop --oknodo --exec /usr/bin/vlc
	eend $?
}

