# main/pingu/pingu.initd
#!/sbin/runscript

description="MultiISP load balancing and failover"
command=/usr/sbin/pingu
command_args="-d"
pidfile=/var/run/pingu/pingu.pid

depend() {
	need net
	after firewall
}
