# testing/rsound/rsound.initd
#!/sbin/openrc-run

name="rsound daemon"
command="/usr/bin/rsd"
command_args="-D $RSOUND_OPTS"
pidfile="/var/run/$SVCNAME.pid"

depend() {
	need net
	use esound
	after firewall
}
