# community/gvfs/gvfs.initd
#!/sbin/openrc-run

command=/usr/libexec/gvfs/gvfsd
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
