# main/nfs-utils/nfsdcld.initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

name="nfsdcld"
description="NFSv4 Client Tracking Daemon"
command=/usr/sbin/nfsdcld

depend() {
    before nfs
    need rpc.pipefs
}

start_pre() {
	checkpath -d -m 0700 -o root:root /var/lib/nfs/nfsdcld
    if ! grep -qs nfsd /proc/filesystems ; then
        modprobe -q nfsd
    fi
}
