# main/nfs-utils/rpc.statd.initd
#!/sbin/openrc-run

[ -e /etc/conf.d/nfs ] && . /etc/conf.d/nfs

name="NFS statd"

command="/sbin/rpc.statd"
command_args="$OPTS_RPC_STATD"
pidfile="/run/rpc.statd.pid"

depend() {
	use ypbind net
	need portmap
	after quota
}

start() {
	# Don't start rpc.statd if already started by someone else...
	killall -q -0 "$command" && return 0

	default_start
}
