# main/nfs-utils/rpc.gssd.initd
#!/sbin/openrc-run

[ -e /etc/conf.d/nfs ] && . /etc/conf.d/nfs

command="/usr/sbin/rpc.gssd"
command_args="$OPTS_RPC_GSSD"

depend() {
	use ypbind net
	need portmap rpc.pipefs
	after quota
}
