# main/nfs-utils/rpc.idmapd.initd
#!/sbin/openrc-run

[ -e /etc/conf.d/nfs ] && . /etc/conf.d/nfs

command="/usr/sbin/rpc.idmapd"
command_args="$OPTS_RPC_IDMAPD"

depend() {
	use ypbind net
	need portmap rpc.pipefs
	after quota
}
