# main/nfs-utils/rpc.svcgssd.initd
#!/sbin/openrc-run

[ -e /etc/conf.d/nfs ] && . /etc/conf.d/nfs

command="/usr/sbin/rpc.svcgssd"
command_args="$OPTS_RPC_SVCGSSD"

depend() {
	use ypbind net
	need portmap rpc.pipefs
	after quota
}
