# community/spotifyd/spotifyd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/spotifyd
command_args="--no-daemon"
command_user="spotifyd:spotifyd"
output_log=/var/log/$RC_SVCNAME.log

depend() {
	need net
}

start_pre() {
	checkpath --directory --owner $command_user --mode 0775 \
		/var/cache/$RC_SVCNAME 
	checkpath --file --owner $command_user --mode 0644 \
		/var/log/$RC_SVCNAME.log
}
