# testing/cortex-tenant/cortex-tenant.initd
#!/sbin/openrc-run

name="cortex-tenant"
description="Prometheus remote write proxy that adds Cortex/Mimir tenant ID based on metric labels"

: ${command_user:="cortex-tenant"}

command="/usr/bin/cortex-tenant"
command_args="-config $CONFIGFILE $command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
