# community/loki/loki.initd
#!/sbin/openrc-run

configfile="${LOKI_CONF:-/etc/loki/loki-local-config.yaml}"

command=/usr/bin/loki
command_args="-config.file=${configfile}
	${LOKI_OPTS}"
command_background=yes
command_user=loki:grafana
make_pidfile=yes
pidfile=/var/run/loki.pid
