# testing/peervpn/peervpn.initd
#!/sbin/openrc-run

# Sample init.d file for alpine linux.

name=peervpn
command="/usr/sbin/$name"
command_args="/etc/$name.conf"
command_background="yes"

# peervpn does not daemonize itself.
start_stop_daemon_args="--stdout=/var/log/$name.log --stderr=/var/log/$name.log"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
