# community/inspircd/inspircd.initd
#!/sbin/openrc-run
name="InspIRCd"
description="InspIRCd IRC server"
supervisor="supervise-daemon"
command=/usr/bin/inspircd
command_user="inspircd:inspircd"
command_args="--nofork --nopid"

depend() {
	need net
	after firewall
}
