# community/udisks2/udisks2.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

description="UDisks2 service"
command="/usr/libexec/udisks2/udisksd"

depend() {
    need dbus
    after dbus
}
