# testing/conserver/conserver.initd
#!/sbin/openrc-run

supervisor="supervise-daemon"
command="/usr/sbin/conserver"
start_stop_daemon_args="--stdout /var/log/${SVCNAME}.log \
	 --stderr /var/log/${SVCNAME}.log"
