# community/dcron/dcron.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/crond"
command_args="$DCRON_OPTS"
pidfile="/var/run/$SVCNAME.pid"
description="Dillon's lightweight cron daemon"

depend() {
	use logger clock hostname
	provide cron
}
