# testing/onedrive/onedrive.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

# onedrive init.d file for alpine linux.
name=onedrive
command=/usr/bin/onedrive
command_args="$onedrive_opts"

depend() {
	need net
	after firewall
}
