# community/conduit/conduit.initd
#!/sbin/openrc-run

command="/usr/bin/conduit"
command_user="conduit:conduit"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

output_log="/var/log/conduit.log"
error_log="/var/log/conduit.log"

depend() {
	need net
}

start_pre() {
	checkpath -d -m 0755 -o "$command_user" /var/lib/conduit
	checkpath -f -m 0644 -o "$command_user" "$output_log"
}
