# testing/aports-mr-to-apk-repo/aports-mr-to-apk-repo.initd
#!/sbin/openrc-run

name="aports-mr-to-apk-repo"
description="Proxy from GitLab CI artifacts to apk repositories"
command="/usr/bin/aports-mr-to-apk-repo"
command_background="yes"
pidfile="/var/run/${RC_SVCNAME}.pid"

depend() {
    need localmount
    after net
}
