# testing/mtail/mtail.initd
#!/sbin/openrc-run

name="mtail"

: "${command_user:="mtail"}"

command="/usr/bin/mtail"
command_args="$ARGS"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

reload() {
	${supervisor} ${RC_SVCNAME} --signal HUP --pidfile "${pidfile}"
}
