# testing/jackal/jackal.initd
#!/sbin/openrc-run

description="Jackal is a server for jackal/XMPP written in Go."
pidfile="/run/jackal.pid"
required_files="/etc/jackal/config.yml"
directory="/var/lib/jackal"
command="/usr/bin/jackal --config /etc/jackal/config.yml"
command_background=true
error_logger="logger -t jackal -p daemon.info"

start_pre() {
	checkpath -d -m 0750 -o jackal:jackal /var/log/jackal
}
