# community/ntpsec/ntpsec.initd
#!/sbin/openrc-run

: ${command_user:=ntp:ntp}

name="ntpsec"
description="NTP reference implementation focused on security"
command="/usr/sbin/ntpd"
command_args="${NTPD_OPTS:--g -N}"
pidfile="/run/$RC_SVCNAME.pid"

# allow setting system time
capabilities="^cap_sys_time"

depend() {
	need net
	provide ntp-client
	use dns
}
