# community/ntpsec/ntpsec.initd
#!/sbin/openrc-run

name="ntpsec"
description="NTP reference implementation focused on security"
command="/usr/sbin/ntpd"
command_args="${NTPD_OPTS:--g -N -u ntp:ntp} -n"
pidfile="/run/$RC_SVCNAME.pid"
command_background="yes"

# allow setting system time
capabilities="^cap_sys_time"

depend() {
	need net
	provide ntp-client
	use dns
}
