# testing/prometheus-smartctl-exporter/smartctl-exporter.initd
#!/sbin/openrc-run
name="prometheus-smartctl-exporter"
command="/usr/bin/smartctl-exporter"
command_args="$ARGS"
command_background="yes"
command_user="prometheus:prometheus"
pidfile="/run/${RC_SVCNAME}.pid"
