# unmaintained/traefik/traefik.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Traefik"
description="Modern reverse proxy and load balancer"

command="/usr/sbin/traefik"
command_args="$traefik_opts"
command_user="traefik:traefik"

description_healthcheck="Check health status by calling /ping endpoint"
healthcheck_timer=30

healthcheck() {
	$command healthcheck >/dev/null
}
