# community/dhcp-helper/dhcp-helper.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
description="DHCP relay agent"

nice="0"
pidfile="/run/dhcp-helper.sd.pid"
supervise_daemon_args=""

# the args are very site-specific, use dhcp-helper --help to find out what you need
ARGS=""

command="/usr/sbin/dhcp-helper"
command_args="-n ${ARGS}"

depend() {
	use net
	after logger firewall
}
