# main/opensmtpd/smtpd.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/$SVCNAME"
command_args="$SMTPD_OPTS"

depend() {
	need net localmount
	after bootmisc firewall
	use logger dns
	provide mta
}

start() {
	ebegin "Starting smtpd ${CONF_MESSAGE}"
	/usr/sbin/smtpd
	eend $?
}

stop() {
       ebegin "Stopping smtpd ${CONF_MESSAGE}"
       /usr/sbin/smtpctl stop &> /dev/null
       eend $?
}