# main/opensmtpd/smtpd.initd
#!/sbin/openrc-run

name="$RC_SVCNAME"
command="/usr/sbin/$RC_SVCNAME"
command_args="-F $SMTPD_OPTS"
command_background=yes
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net localmount
	after bootmisc firewall
	use logger dns
	provide mta
}
