# testing/nscd/nscd.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/nscd"
command_args="${NSCD_OPTS:--d -p /var/run/$SVCNAME.pid}"
pidfile="/var/run/$SVCNAME.pid"

start_pre() {
        mkdir -p /var/run/nscd
}
