# testing/diod/diod.initd
#!/sbin/runscript

start() {
	ebegin
		start-stop-daemon --start -p /var/run/diod.pid -mbx /usr/bin/diod -- -f
	eend $?
}
