# community/sniproxy/sniproxy.initd
#!/sbin/openrc-run

pidfile="/var/run/${SVCNAME}/${SVCNAME}.pid"
command="/usr/sbin/${SVCNAME}"
: ${conf_file:="/etc/sniproxy/sniproxy.conf"}

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting $SVCNAME"
	start-stop-daemon --start \
		--quiet --pidfile $pidfile --exec $command \
		-- -c "$conf_file" "$@"
	eend
}

