# testing/moosefs/moosefs-master.initd
#!/sbin/openrc-run

command=/usr/sbin/mfsmaster

depend() {
	need syslog
	after net
}
