# community/sensorfw/sensorfw.initd
#!/sbin/openrc-run

name="SensorFW"
description="Sensor Framework Qt5"

command=/usr/sbin/sensorfwd

depends() {
	need udev
}
