# community/sensorfw/sensorfw.initd
#!/sbin/openrc-run

name="SensorFW"
description="Sensor Framework Qt5"

command=/usr/sbin/sensorfwd
command_args="-c=/etc/sensorfw/primaryuse.conf --no-magnetometer-bg-calibration"

depend() {
	need udev
}
