# community/connman/connman.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Network Management Daemon"
description="Daemon for managing internet connections"

: ${cfgfile:=/etc/connman/main.conf}

command=/usr/sbin/connmand
command_args="$command_args -c $cfgfile"
command_args_foreground="--nodaemon"

depend() {
	need dbus
	provide net
}
