# community/connman/connman.initd
#!/sbin/openrc-run

name="ConnMan"
description="Daemon for managing internet connections"

: ${cfgfile:=/etc/connman/main.conf}

command=/usr/sbin/connmand
command_args="$command_args -c $cfgfile --nodaemon"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
required_files="$cfgfile"

depend() {
	need dbus
	provide net
}
