# community/rng-tools/rngd.initd
#!/sbin/openrc-run

description="Random number generator daemon"

command="/usr/sbin/rngd"
command_args="-f ${command_args:-$RNGD_OPTS}"
command_background="yes"

start_stop_daemon_args="--wait 5"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need localmount
	after urandom
	provide entropy
}
