# community/rng-tools/rngd.initd
#!/sbin/openrc-run

description="Random number generator daemon"

pidfile="/run/$RC_SVCNAME.pid"
command="/usr/sbin/rngd"
command_args="--quiet --background --pid-file $pidfile $RNGD_OPTS"
start_stop_daemon_args="--wait 5"

depend() {
	need localmount
	after urandom
	provide entropy
}
