# testing/piping-server/piping-server.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Piping Server"
description="Infinitely transfer data between devices over HTTP"

command="/usr/bin/piping-server"
command_user="nobody:nobody"
directory="/var/empty"

: ${error_log:="/var/log/piping-server.log"}

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f -o $command_user "$error_log"
}
