# testing/wroomd/wroomd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="wroomd"
description="Modem management daemon for Qualcomm platforms"

command="/usr/bin/wroomd"

depend() {
	need dbus
}
