# community/kodi/kodi.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Kodi"
description="A software media player and entertainment hub for digital media"

command=/usr/lib/kodi/kodi-gbm
command_args="--standalone"
command_background=true
pidfile="/var/run/kodi.pid"
command_user="kodi:kodi"
