# community/coventry/coventry.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name=coventry
command="/usr/sbin/coventry"
pidfile="/run/coventry.pid"
extra_started_commands="reload"

depend() {
        use logger dns
        need net
        after firewall
}

reload() {
  ebegin "Reloading ${RC_SVCNAME}"
  start-stop-daemon --signal HUP --pidfile "${pidfile}"
  eend $?
}

