# community/user-managerd/user-managerd.initd
#!/sbin/openrc-run

name=user-managerd
description="Sailfish user manager service"
command="/usr/bin/user-managerd"
command_background=yes

depend() {
	need dbus
}
