# community/polkit-elogind/polkit.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Polkit System Daemon"
description="Provide System authentication via PolicyKit1 D-Bus service"

command=/usr/lib/polkit-1/polkitd
command_args="${POLKIT_OPTS:---no-debug}"

depend() {
	need dbus
}
