# main/cacti/cacti.crontab
*/5 * * * * php /usr/share/webapps/cacti/poller.php >/dev/null 2>&1
