# main/nrpe/nrpe.initd
#!/sbin/openrc-run

: ${user:="nagios"}
: ${cfgfile:="/etc/nrpe.cfg"}
: ${nrpe_options:=" -c $cfgfile"}  # for backward compatibility

command="/usr/bin/nrpe"
command_args="-f $nrpe_options ${command_args:-}"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="
	--user $user
	${group:+--group $group}"

depend() {
	need net
	after firewall
}
