# community/syncthing/syncthing.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
supervisor=supervise-daemon
command=/usr/bin/syncthing
command_args="${SYNCTHING_ARGS:---no-browser --no-restart}"
command_user="${SYNCTHING_USER:-syncthing}"
output_logger="logger -t '${RC_SVCNAME}' -p daemon.info"
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
        use logger dns net
        after firewall
}
