# community/syncthing/syncthing-user.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
command=/usr/bin/syncthing
: "${command_args:=--no-browser --no-restart}"
output_logger="logger -t '${RC_SVCNAME}' -p daemon.info"
