# unmaintained/kde-workspace/kdm.initd
#!/sbin/runscript

depends() {
	needs localmount dbus
}

description="K Display Manager"
command="/usr/bin/kdm -nodaemon"
pidfile="/var/run/kdm.pid"

