# testing/ajaxterm/ajaxterm.initd
#!/sbin/openrc-run

depend()
{
    need net
}

start()
{
    ebegin "Starting AjaxTerm on port $PORT"
    start-stop-daemon --start --pidfile $PIDFILE --exec $DAEMON -- --daemon --port=$PORT --uid=nobody
    eend $?
}

stop()
{
    ebegin "Stopping AjaxTerm"
    start-stop-daemon  --stop --pidfile $PIDFILE 
    rm -f $PIDFILE 
    eend $?
}
