# testing/mcron/mcron.initd
#!/sbin/openrc-run

name="$SVCNAME"
pidfile="/run/$SVCNAME.pid"
command="/usr/bin/$SVCNAME"
command_args="$MCRON_OPTS"
command_background=true
error_logger="logger -t mcron"

depend() {
	use logger
}
