# testing/airsonic/airsonic.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="airsonic"
command="/usr/bin/java"
command_args="-jar airsonic.war"
command_user="airsonic:airsonic"
directory="/var/airsonic"
command_background=true
pidfile="/run/airsonic.pid"

depend() {
	need net localmount
	after firewall
}
