# testing/frp/frps.initd
#!/sbin/openrc-run

name=frps
description="frp server"
command="/usr/bin/$name"
supervisor=supervise-daemon

: ${cfgfile:=/etc/frp/$RC_SVCNAME.ini}
command_args="-c $cfgfile"
pidfile="/run/$RC_SVCNAME.pid"

: ${respawn_delay:=5}
: ${respawn_max:=0}

depend() {
	need net
	after firewall
}
