# testing/dma/dma.cron
#!/bin/sh
# Flush the DMA queue automatically.

/usr/sbin/dma -q1
exit 0
