# testing/speedtest_exporter/speedtest_exporter.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=speedtest_exporter
description="Prometheus exporter for Speedtest metrics"

command=/usr/bin/speedtest_exporter
command_args="$speedtest_exporter_opts"
command_user="speedtest-exporter:speedtest-exporter"

depend() {
	need net
	after firewall
}
