# testing/taskcafe/taskcafe.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="taskcafe"
command="/usr/bin/taskcafe"
command_args="--config $taskcafe_configfile web"
command_user="taskcafe:taskcafe"
command_background=true
pidfile="/run/taskcafe.pid"

depend() {
	need net localmount
	after firewall
}
