# testing/influxdb/influxdb.initd
#!/sbin/openrc-run

pidfile=/run/influxdb.pid
start_stop_daemon_args="--user influxdb --stderr /var/log/influxdb/influxdb.log"
command=/usr/sbin/influxd
command_args="${EXTRA_OPTS}"
command_background="yes"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -o influxdb:influxdb -m755 /var/lib/influxdb /var/log/influxdb
}
