# testing/firehol/firehol.initd
#!/sbin/openrc-run

name=firehol
daemon=/usr/sbin/$name

depend() {
	need net
}

start() {
	ebegin "Starting ${name}"
	firehol start
	eend $?
}

stop() {
	ebegin "Stopping ${name}"
	firehol stop
	eend $?
}

try() {
	ebegin "Testing ${name} configuration"
	firehol try
	eend $?
}

panic() {
	ebegin "Putting ${name} in PANIC mode"
	firehol panic "$@"
	eend $?
}

