# community/waybar/waybar.initd
#!/sbin/openrc-run

command=/usr/bin/waybar
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
	use wireplumber
}
