# testing/hyprpolkitagent/hyprpolkitagent.initd
#!/sbin/openrc-run

description="A polkit authentication agent written in QT/QML"

command="/usr/libexec/hyprpolkitagent"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
