# main/wpa_supplicant/wpa_cli.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="WPA Command Line Client"
description="Text-based frontend for interacting with WPA Supplicant"

command=/sbin/wpa_cli
command_args="${WPACLI_OPTS} >/dev/null"

depend() {
	need wpa_supplicant
}
