# testing/eg25-manager/eg25-manager.initd
#!/sbin/openrc-run
supervisor="supervise-daemon"
description="Starting EG25 Manager"
command="/usr/bin/eg25manager"

depend()
{
        after modemmanager
}
