# main/open-lldp/lldpad.initd
#!/sbin/openrc-run

name="lldpad"
command="/usr/sbin/lldpad"
command_args="$command_args -p"
command_background=true
pidfile=/run/$RC_SVCNAME.pid

depend() {
	need net
	after firewall
}
