# testing/openlldp/openlldp.initd
#!/sbin/runscript

NAME="lldpd"
DAEMON="/usr/sbin/$NAME"
pidfile=/var/run/$NAME.pid

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting ${NAME}"
	start-stop-daemon --start --quiet --pidfile ${pidfile} \
		--exec ${DAEMON} -- ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"
	start-stop-daemon --stop --exec ${DAEMON}
	eend $?
}

