# community/zabbix/zabbix-agentd.initd
#!/sbin/openrc-run

command=/usr/sbin/zabbix_agentd
command_args_foreground="-f"
supervisor="supervise-daemon"

start_pre() {
        checkpath -d -m 0775 -o zabbix:zabbix /var/run/zabbix
}
