# community/zabbix/zabbix-agentd.initd
#!/sbin/openrc-run

command=/usr/sbin/zabbix_agentd
command_args_foreground="-f"
supervisor="supervise-daemon"

