# community/zabbix/zabbix-server.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Monitoring Server"
description="Daemon for the Zabbix Monitoring Server"

command=/usr/bin/zabbix_server
command_args_foreground="-f"

command_user=zabbix:zabbix

depend() {
	need net
	use mysql postgresql
	after firewall
}

start_pre() {
	checkpath -d -m 0775 -o $command_user /var/run/zabbix
}
