# community/zabbix/zabbix-agent2.initd
#!/sbin/openrc-run

command=/usr/sbin/zabbix_agent2
supervisor="supervise-daemon"
command_user="zabbix:zabbix"

start_pre() {
    checkpath --directory --owner zabbix:zabbix /var/log/zabbix
    checkpath --directory --owner zabbix:zabbix /run/zabbix
}
