# community/zabbix/zabbix-agent2.initd
#!/sbin/openrc-run

command=/usr/sbin/zabbix_agent2
supervisor="supervise-daemon"

pre_start() {
    checkpath --directory --owner zabbix:zabbix /var/log/zabbix
}
