# testing/pimd-dense/pimd-dense.initd
#!/sbin/openrc-run

name=pimdd
command="/usr/sbin/$name"
command_args="$pimdd_opts -n"
command_background="yes"

pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
