# testing/wok/wokd.initd
#!/sbin/openrc-run
# (c) 2016 Alpine Linux

name=wokd
command=/usr/bin/$name
command_background=YES
pidfile=/run/$name.pid

depend() {
        use logger dns
        need net
        after libvirtd firewall
}
