# community/liquidctl/liquidctl.initd
#!/sbin/openrc-run

description="Liquid Cooler Control"
command=/usr/bin/liquidctl
: ${command_user:=liquidctl:liquidctl}

start_pre() {
	if [ -z "${command_args}" ]; then
		eerror "command_args not specified in /etc/conf.d/liquidctl"
		return 1
	fi
}

no_new_privs="yes"
