# testing/fluent-bit/fluent-bit.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

description="Fast and Lightweight Log/Data Forwarder [fluent-bit]"

nice=19
pidfile="/run/$SVCNAME.sd.pid"
supervise_daemon_args="-p $pidfile -N $nice"
command=/usr/bin/td-agent-bit
command_args=${TD_OPTS}

depends() {
	use net
}
