# testing/fluent-bit/fluent-bit.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=fluent-bit
command="/usr/bin/fluent-bit"
command_args="$fluentbit_opts"

depend() {
	need net
	after firewall
}
