# main/taskd/taskd.initd
#!/sbin/openrc-run

description="Secure server providing multi-user, multi-client access to task data"
command=/usr/bin/taskd
command_args="server --data /var/lib/taskd"

start_stop_daemon_args="--user taskd:taskd"

depend() {
	need net
	after firewall
}
