# community/cosmic-greeter/cosmic-greeter-daemon.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Greeter Daemon"
description="COSMIC Greeter Daemon"

command=/usr/bin/cosmic-greeter-daemon

depend() {
	need dbus
	need seatd
}
