# community/cosmic-greeter/cosmic-greeter.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Greeter"
description="COSMIC Greeter"

command=/usr/sbin/greetd
command_args="--config /etc/greetd/cosmic-greeter.toml"

depend() {
	provide display-manager
	need cosmic-greeter-daemon
	need dbus
	need localmount
}
