# testing/mediaproxy/mediaproxy-relay.initd
#!/sbin/runscript

daemon=/usr/bin/media-relay
pidfile=/var/run/mediaproxy/relay.pid

depend() {
	need net
	after firewall
	after mediaproxy-dispatcher
}

start() {
	ebegin "Starting Mediaproxy Relay"
	start-stop-daemon --start --quiet --pidfile $pidfile --exec $daemon 
	eend $?
}

stop() {
	ebegin "Stopping Mediaproxy Relay"
	start-stop-daemon --stop --quiet --pidfile $pidfile
	eend $?
}
