# community/bareos/bareos-sd.initd
#!/sbin/openrc-run

name="Bareos Storage"
command=/usr/sbin/bareos-sd
pidfile=/run/$RC_SVCNAME.pid

depend() {
	need net
	use dns
}
