# community/bareos/bareos-fd.initd
#!/sbin/openrc-run

name="Bareos File"
command=/usr/sbin/bareos-fd
pidfile=/run/$RC_SVCNAME.pid

depend() {
	need net
	use dns
}
