# community/bareos/bareos-fd.initd
#!/sbin/openrc-run

name="Bareos File daemon"
pidfile=/run/bareos/bareos-fd.9102.pid
command=/usr/sbin/bareos-fd
command_args="-u root -g bareos"

depend() {
	need net
	use dns
}

start_pre() {
	checkpath -dm750 -o bareos "${pidfile%/*}"
}