# community/bareos/bareos-dir.initd
#!/sbin/openrc-run

name="Bareos Director"
command=/usr/sbin/bareos-dir
pidfile=/run/$RC_SVCNAME.pid

depend() {
	need net
	after firewall
	use dns bareos-fd bareos-sd postgresql mariadb
}
