# testing/fastd/fastd.initd
#!/sbin/openrc-run

VPN=${SVCNAME#*.}
if [ ${SVCNAME} != "fastd" ]; then
  pidfile="/var/run/fastd.${VPN}.pid"
else
  pidfile="/var/run/fastd.pid"
fi

command=/usr/bin/fastd
command_args="--daemon --pid-file ${pidfile} --config /etc/fastd/${VPN}.conf"
if [ ${SVCNAME} == "fastd" ]; then
  name="fastd"
else
  name="fastd ${VPN}"
fi
description="fastd VPN daemon"


