# testing/ruuvi-prometheus/ruuvi-prometheus.initd
#!/sbin/openrc-run

pidfile="/var/run/${RC_SVCNAME}.pid"
command="/usr/bin/${RC_SVCNAME}"
command_background=true

depend() {
	after coldplug
}
