# testing/babeld/babeld.initd
#!/sbin/openrc-run

command="/usr/bin/babeld"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-D -I $pidfile $ARGS"

depend() {
	need net
	after firewall
}

start_pre() {
	if [ -z "$ARGS" ] && [ ! -s "/etc/babeld.conf" ]; then
		echo "ERROR: No configuration found in /etc/conf.d/babeld or /etc/babeld.conf"
		return 1
	fi
}
