# testing/nebula/nebula.initd
#!/sbin/openrc-run
name=$RC_SVCNAME
description="Scalable overlay networking tool"
cfgfile="/etc/$RC_SVCNAME/config.yml"
command="/usr/sbin/$RC_SVCNAME"
command_args="-config $cfgfile"
command_user="$RC_SVCNAME"
supervisor="supervise-daemon"

depend() {
	need net
	use logger dns
	after firewall
	provide nebula
}

start_pre() {
	$command -config $cfgfile -test
}
