# main/openrc/hostname.initd
#!/sbin/openrc-run

# Copyright (c) Natanael Copa
# This code is licensed under BSD-2-Clause

description="Sets the hostname of the machine."

depend() {
	keyword -prefix -lxc -docker
}

start() {
	if [ -s /etc/hostname ] ; then
		opts="-F /etc/hostname"
	else
		opts="${hostname:-localhost}"
	fi
	ebegin "Setting hostname"
	hostname $opts
	eend $?
}
