# main/openrc/keymaps.initd
#!/sbin/runscript

description="Applies a keymap for the consoles."

depend()
{
	need localmount
	keyword noopenvz noprefix nouml novserver noxenu
}

start() {
	[ -z "$KEYMAP" ] && return
	ebegin "Setting keymap"
	zcat "$KEYMAP" | loadkmap
	eend $?
}

stop() {
	return
}
