# main/openrc/machine-id.initd
#!/sbin/openrc-run

# Copyright (c) Natanael Copa
# This code is licensed under BSD-2-Clause

description="Generate machine-id if needed"

depend() {
	need root dev
}

start() {
	if [ -s /etc/machine-id ] ; then
		return 0
	fi
	ebegin "Generating machine-id"
	dd if=/dev/urandom status=none bs=16 count=1 \
		| md5sum | cut -d' ' -f1 > /etc/machine-id
	eend $?
}
