# main/openrc/modules.initd
#!/sbin/runscript

description="Loads a user defined list of kernel modules."

depend()
{
	before hwclock hwdrivers
	keyword noopenvz noprefix novserver
}


start() {
	if [ -f /etc/modules ] ; then
		ebegin "Loading modules"
		sed 's/\#.*//g' < /etc/modules |
		while read module args
		do
			modprobe -q $module $args
		done
		eend $?
	fi
}

