# testing/acf2/acf2.initd
#!/sbin/runscript

# init.d file for acf2
# Copyright (c) 2013 Kaarle Ritvanen

name=acf2
command=/usr/sbin/$name
pidfile=/var/run/$name.pid

depend() {
	need net
	after firewall
}

stop() {
	ebegin "Stopping $name"
	start-stop-daemon --stop --pidfile $pidfile --signal INT
	eend $?
}
