# community/usbguard/usbguard.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="USBGuard"
description="Software framework for implementing USB device authorization policies"
command="/usr/sbin/usbguard-daemon"
command_args="-s"

depend() {
	use logger
}
