# community/usbguard/usbguard-dbus.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="USBGuard D-Bus interface service"
description="Software framework for implementing USB device authorization policies"
command="/usr/sbin/usbguard-daemon"
command_args="--system"

depend() {
	need dbus usbguard
}
