# testing/tvheadend-git/tvheadend-git.initd
#!/sbin/openrc-run

command="/usr/bin/tvheadend"
command_args="-f -c $TVHCONF -u tvheadend -g video $TVHARGS"
retry="TERM/30/KILL/5"
command_background=true
pidfile=$TVHPID
name="tvheadend-git"
required_dirs=$TVHCONF

depend() {
	use logger dns
	need net
	after firewall
}
