# main/busybox-initscripts/acpid.initd
#!/sbin/runscript

depend() {
	need dev localmount
	after hwdrivers modules
}

start() {
	ebegin "Starting busybox acpid"
	start-stop-daemon --start --exec /sbin/acpid \
		-- $ACPID_OPTS
	eend $?
}

stop () {
	ebegin "Stopping busybox acpid"
	start-stop-daemon --stop --exec /sbin/acpid
	eend $?
}

