# main/busybox-initscripts/acpid.initd
#!/sbin/runscript

name="busybox acpid"
command="/sbin/acpid"
command_args="$ACPID_OPTS"
pidfile="/var/run/acpid.pid"

depend() {
	need dev localmount
	after hwdrivers modules
}

