# main/busybox-initscripts/cron.initd
#!/sbin/runscript

depend() {
	need localmount
	need logger
}

start() {
	ebegin "Starting busybox cron"
	start-stop-daemon --start --exec /usr/sbin/crond -- $CRON_OPTS
	eend $?
}

stop () {
	ebegin "Stopping busybox cron"
	start-stop-daemon --stop --exec /usr/sbin/crond
	eend $?
}

