# main/busybox-initscripts/dnsd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="Small static dns server daemon"

command="/usr/sbin/$RC_SVCNAME"
command_args="$DNSD_OPTS"

depend() {
	need net
	after firewall
}

