# main/busybox-initscripts/dnsd.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting dnsd"
	start-stop-daemon --start --exec /usr/sbin/dnsd -- -d $DNSD_OPTS
	eend $?
}

stop () {
	ebegin "Stopping dnsd"
	start-stop-daemon --stop --exec /usr/sbin/dnsd
	eend $?
}
