# main/busybox-initscripts/httpd.initd
#!/sbin/runscript

DAEMON=/usr/sbin/httpd

depend() {
	need net localmount
	after firewall
}

start() {
	ebegin "Starting busybox httpd"
	start-stop-daemon --start --exec $DAEMON -- $HTTPD_OPTS
	eend $?
}

stop () {
	ebegin "Stopping busybox httpd"
	start-stop-daemon --stop --exec $DAEMON
	eend $?
}

