# main/busybox-initscripts/inetd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="Listen for network connections and launch programs"

command="/usr/sbin/$RC_SVCNAME"
command_args="$INETD_OPTS"
command_args_foreground="-f"

depend() {
	need net localmount
	after firewall
}

