# main/busybox-initscripts/inetd.initd
#!/sbin/openrc-run

name="busybox $SVCNAME"
command="/usr/sbin/$SVCNAME"
command_args="-f $INETD_OPTS"
pidfile="/var/run/$SVCNAME.pid"
command_background=true

depend() {
	need net localmount
	after firewall
}

