# main/busybox-initscripts/inetd.initd
#!/sbin/runscript

depend() {
	need net localmount
	after firewall
}

start() {
	ebegin "Starting busybox inetd"
	start-stop-daemon --start --exec /usr/sbin/inetd -- $INETD_OPTS
	eend $?
}

stop () {
	ebegin "Stopping busybox inetd"
	start-stop-daemon --stop --exec /usr/sbin/inetd
	eend $?
}

