# main/busybox-initscripts/klogd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="Daemon for intercepting and logging Linux Kernel messages"

command="/sbin/$RC_SVCNAME"
command_args="${KLOGD_OPTS}"
command_args_foreground="-n"

depend() {
	need clock hostname localmount syslog
	before net
	keyword -vserver -lxc
}
