# main/busybox-initscripts/klogd.initd
#!/sbin/runscript

depend() {
	need clock hostname localmount
	before net
	provide logger
	keyword novserver
}

start() {
	ebegin "Starting busybox kernel logging"
	start-stop-daemon --start --exec /sbin/klogd -- ${KLOGD_OPTS}
	eend $?
}

stop () {
	ebegin "Stopping busybox kernel logging"
	start-stop-daemon --stop --exec /sbin/klogd
	eend $?
}

