# main/busybox-initscripts/ntpd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="Network Time Protocol Daemon Server"

command="/usr/sbin/$RC_SVCNAME"
command_args="${NTPD_OPTS:--N -p pool.ntp.org}"
command_args_foreground="-n"

depend() {
	need net
	provide ntp-client
	use dns
}
