# main/busybox-initscripts/syslog.initd
#!/sbin/runscript

depend() {
	need clock hostname klogd localmount
	provide logger
}

start() {
	ebegin "Starting busybox system logging"
	start-stop-daemon --start --exec /sbin/syslogd -- ${SYSLOGD_OPTS}
	eend $?
}

stop () {
	ebegin "Stopping busybox system logging"
	start-stop-daemon --stop --exec /sbin/syslogd
	eend $?
}

