# main/busybox-initscripts/syslog.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox syslog"
description="Message logging system"

command="/sbin/syslogd"
command_args="${SYSLOGD_OPTS}"
command_args_foreground="-n"

supervise_daemon_args="-g wheel -k 027"

depend() {
	need clock hostname localmount
	provide logger
}
