# main/busybox-initscripts/syslog.initd
#!/sbin/runscript

name="busybox syslog"
command="/sbin/syslogd"
command_args="${SYSLOGD_OPTS}"
pidfile="/var/run/syslogd.pid"

depend() {
	need clock hostname klogd localmount
	provide logger
}
