# main/busybox-initscripts/udhcpd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="DHCP server"

command="/usr/sbin/$RC_SVCNAME"
command_args="$UDHCPD_OPTS"
command_args_foreground="-f"

depend() {
	need net
	after firewall
}
