# main/busybox-initscripts/watchdog.initd
#!/sbin/runscript

depend() {
	need dev
	after hwdrivers
}

start() {
	if ! [ -n "$WATCHDOG_DEV" ]; then
		eerror "WATCHDOG_DEV is not set"
		return 1
	fi

	ebegin "Starting busybox watchdog"
	start-stop-daemon --start --exec /sbin/watchdog \
		-- $WATCHDOG_OPTS $WATCHDOG_DEV
	eend $?
}

stop () {
	ebegin "Stopping busybox watchdog"
	start-stop-daemon --stop --exec /sbin/watchdog
	eend $?
}

