# main/busybox-initscripts/watchdog.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="busybox $RC_SVCNAME"
description="software watchdog daemon"

command="/sbin/$RC_SVCNAME"
command_args="${WATCHDOG_OPTS} ${WATCHDOG_DEV}"
command_args_foreground="-F"

depend() {
	need dev
	after hwdrivers
}

start_pre() {
	if ! [ -n "$WATCHDOG_DEV" ]; then
		eerror "WATCHDOG_DEV is not set"
		return 1
	fi
}
