# community/snowflake/snowflake-proxy.initd
#!/sbin/openrc-run
name="snowflake-proxy"
description="pluggable transport for Tor using WebRTC"
supervisor=supervise-daemon
command=/usr/bin/snowflake-proxy
command_args="-log $output_log $snowflake_args"
command_user="snowflake:snowflake"

start_pre() {
	[ -n "$output_log" ] && checkpath -f "$output_log" \
		-m 644 -o snowflake:snowflake
}
