# testing/cloudflared/cloudflared.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=cloudflared
description="cloudflared - Argo Tunnel Client"

command=/usr/sbin/cloudflared
command_args="${cloudflared_opts}"
command_user="cloudflared:cloudflared"
pidfile=/run/${RC_SVCNAME}.pid
