# testing/cloudflared/cloudflared.initd
#!/sbin/openrc-run

name=cloudflared
command=/usr/bin/cloudflared
command_user=cloudflared:cloudflared
command_background="yes"
pidfile=/run/${RC_SVCNAME}.pid

depend() {
	need net
	after firewall
}
