# community/umtprd/umtprd.initd
#!/sbin/openrc-run

name=usb-moded
description="usb-moded USB gadget controller"
command="/usr/sbin/usb_moded"
command_background=yes

depend() {
	need dbus
}
