# testing/microsocks/microsocks.initd
#!/sbin/openrc-run

name=microsocks
description="multithreaded, small, efficient SOCKS5 server"
command="/usr/bin/$name"
supervisor=supervise-daemon

pidfile="/run/$RC_SVCNAME.pid"

: ${respawn_delay:=5}
: ${respawn_max:=0}

depend() {
	need net
	after firewall
}
