# testing/sddm/sddm.openrc
#!/sbin/openrc-run

depends() {
	needs localmount
}

description="SDDM Display Manager"
command="/usr/bin/sddm"
command_args="$sddm_opts"
command_background="true"
pidfile="/var/run/sddm.pid"
