# community/ksmbd-tools/ksmbd.initd
#!/sbin/openrc-run

command=/usr/sbin/ksmbd.mountd
command_args="$KSMBD_MOUNTD_OPTS"
pidfile=/var/run/$SVCNAME.lock

depend() {
	need net localmount
	after firewall
}

start_pre() {
	if [ ! -d /sys/class/ksmbd-control ]; then
		if ! modprobe ksmbd; then
			return 1
		fi
	fi
}

stop() {
	/usr/sbin/ksmbd.control -s
}
