# community/ksmbd-tools/ksmbd.initd
#!/sbin/openrc-run

command=/usr/sbin/ksmbd.mountd
command_args="$KSMBD_MOUNTD_OPTS"
pidfile=/var/run/$SVCNAME.pid

depend() {
	need net localmount
	after firewall
}
stop() {
	/usr/sbin/ksmbd.control -s
}
