# unmaintained/tcpproxy/tcpproxy.initd
#!/sbin/openrc-run

DAEMON=/usr/sbin/tcpproxy
start() {
	ebegin "Starting tcpproxy"
	start-stop-daemon -S -x $DAEMON -- $TCPPROXY_OPTS
	eend $?
}

stop() {
	ebegin "Stopping tcpproxy"
	start-stop-daemon -K -x $DAEMON
	eend $?
}
