# testing/spice-vdagent/spice-vdagentd.initd
#!/sbin/openrc-run

command=/usr/sbin/spice-vdagentd
command_args=""
pidfile="/run/spice-vdagentd/spice-vdagentd.pid"

depend() {
	need dbus
}

start_pre() {
	checkpath --directory --mode 0755 --owner root:root /run/spice-vdagentd
}
