# testing/zigbee2mqtt/zigbee2mqtt.initd
#!/sbin/openrc-run

name="Zigbee2MQTT"

: ${command_user:="zigbee2mqtt"}

command="/usr/bin/zigbee2mqtt"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

export ZIGBEE2MQTT_DATA="${datadir:="/var/lib/zigbee2mqtt"}"

depend() {
	need dev localmount
	use mosquitto
}
