# testing/mlmym/mlmym.initd
#!/sbin/openrc-run

name=mlmym
description="Mlmym is a Lemmy frontend offering a familiar desktop experience"

export LEMMY_DOMAIN

command="/usr/bin/mlmym"
command_args="${command_args:---addr 0.0.0.0:8008}"
command_background="yes"
command_user="mlmym:mlmym"
pidfile="/run/${RC_SVCNAME}.pid"
directory="/usr/share/webapps/mlmym"

depend() {
	need net localmount
	after firewall
}
