# community/do-agent/do-agent.initd
#!/sbin/openrc-run
name="DigitalOcean Monitoring Agent"
command=/usr/sbin/do-agent
command_user=do-agent
command_background=true
pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="--exec $RC_SVCNAME"

depend() {
	need net
}
