# community/gnome-keyring/gnome-keyring.initd
#!/sbin/openrc-run

description="GNOME keyring daemon"

command="/usr/bin/gnome-keyring-daemon"
command_args="--foreground --components=$KEYRING_COMPONENTS --control-directory=$KEYRING_DIRECTORY"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
