# main/in-sync/in-sync.initd
#!/sbin/openrc-run

# inotify-triggered file synchronization
# Copyright (c) 2020 Kaarle Ritvanen

name=in-sync
command=/usr/sbin/in-sync
command_background=1
pidfile=/var/run/$name.pid

depend() {
	need net
	after firewall
}
