# testing/mdnsd/mdnsd.initd
#!/sbin/openrc-run

description="embeddable multicast DNS daemon"
command="/usr/sbin/mdnsd"
command_args="${MDNSD_OPTS}"
command_background="yes"
command_user="daemon:daemon"

pidfile="/run/${RC_SVCNAME}.pid"
extra_started_commands="reload"

depend() {
	need net
	after firewall
}

reload() {
	ebegin "Reloading ${RC_SVCNAME}"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
