# testing/zrepl/zrepl.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
cfgfile="/etc/$RC_SVCNAME/$RC_SVCNAME.yml"
command="/usr/bin/zrepl"
command_args=" --config /etc/zrepl/zrepl.yml daemon"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"
command_background="yes"

start_pre() {
	checkpath --directory --mode 0700 /run/$RC_SVCNAME
}
