# testing/hexagonrpcd/hexagonrpcd-sdsp.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Hexagonrpcd SDSP"
description="Daemon to support Qualcomm Hexagon SDSP virtual filesystem"

command=/usr/bin/hexagonrpcd
command_args="-f /dev/fastrpc-sdsp -s"
command_user=fastrpc:fastrpc
