# testing/hexagonrpcd/hexagonrpcd-adsp-sensorspd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Hexagonrpcd ADSP SensorPD"
description="Daemon to support Qualcomm Hexagon ADSP virtual filesystem for SensorPD"

command=/usr/bin/hexagonrpcd
command_args="-f /dev/fastrpc-adsp -s"
command_user=fastrpc:fastrpc
