# non-free/3dm2/3dm2.initd
#!/sbin/runscript

NAME=3dm2
DAEMON=/usr/sbin/$NAME

depend() {
        need net
	after firewall
}

start() {
	ebegin "Starting ${NAME}"
	start-stop-daemon --start --quiet --background \
		--exec ${DAEMON} -- ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"
	start-stop-daemon --stop --quiet \
		--exec ${DAEMON} \
	eend $?
}
