# testing/unrealircd/unrealircd.initd
#!/sbin/openrc-run
name="Unrealircd"
supervisor=supervise-daemon
command_user=unrealircd
command=/usr/bin/unrealircd
command_args="-F"
directory=/var/lib/unrealircd

depend() {
        use net
}
