# testing/unrealircd/unrealircd.initd
#!/sbin/openrc-run
name="Unrealircd"
command_user=unrealircd
command=/usr/bin/unrealircd
command_args="-F $command_args"
directory=/var/lib/unrealircd

depend() {
        use net
}
