# testing/polipo/polipo.initd
#!/sbin/openrc-run
depend() {
	use net
}

start() {
 	ebegin "Starting polipo"
	start-stop-daemon --start --background --pidfile /var/run/polipo.pid --make-pidfile --exec /usr/bin/polipo
	eend $?
}

stop() {
 	ebegin "Stopping polipo"
	start-stop-daemon --stop --pidfile /var/run/polipo.pid
	eend $?
}
