# main/charybdis/charybdis.initd
#!/sbin/openrc-run

command="/usr/bin/charybdis"
command_user="${CHARYBDIS_USER:-ircd}"

depend() {
	use net
	provide ircd
}
