# testing/kopano-core/kopano-ical.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-ical"

depend() {
	need kopano-server
}
