# testing/kopano-core/kopano-dagent.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-dagent"

depend() {
	need net kopano-server
}

start_pre() {
	checkpath -d  -m 750 -o kopano-core:kopano /var/lib/kopano/dagent
}
