# testing/kopano-core/kopano-monitor.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-monitor"

depend() {
	need kopano-server
}
