# testing/kopano-core/kopano-spooler.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-spooler"

depend() {
	need kopano-server
}
