# testing/kopano-core/kopano-spamd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-spamd"

depend() {
	need kopano-server
}
