# testing/kopano-core/kopano-server.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
command="/usr/bin/kopano-server"

depend() {
	need mariadb
}

start_pre() {
	checkpath -d -m 750 -o kopano-core:kopano /run/kopano
}
