# community/slim/slim.initd
#!/sbin/openrc-run

depend() {
	need localmount dbus
}

start() {
	ebegin "Starting slim"
	start-stop-daemon --start --quiet --exec /usr/bin/slim -- -d
	eend $?
}

stop() {
	ebegin "Stopping slim"
	start-stop-daemon --stop --quiet --exec /usr/bin/slim
	eend $?
}

