# testing/hddfancontrol/hddfancontrol.initd
#!/sbin/openrc-run

name=hddfancontrol
description="Regulate fan speed according to hard drive temperature"

command=/usr/bin/hddfancontrol
command_args="-v $HDDFANCONTROL_LOG_LEVEL daemon $HDDFANCONTROL_DAEMON_ARGS"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	use lm_sensors
}
