# main/rsync/rsyncd.initd
#!/sbin/openrc-run

cfgfile="/etc/rsyncd.conf"
pidfile="$(grep "^pid file" $cfgfile | cut -d= -f2)"
command="/usr/bin/rsync"
command_args="--daemon $command_args $RSYNC_OPTS"
required_files="$cfgfile"

depend() {
	use net
	after firewall
}
