# testing/grlx/grlx-sprout.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

description="Grlx Sprout"
command="/usr/bin/grlx-sprout"

depend() {
	need net
}
