# testing/grlx/grlx-farmer.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
description="Grlx Farmer"
command="/usr/bin/grlx-farmer"

depend() {
	need net
}
