# community/icingaweb2-module-director/icinga-director.initd
#!/sbin/openrc-run
# (c) 2021 Francesco Colista <fcolista@alpinelinux.org>

name=$RC_SVCNAME
command="/usr/bin/icingacli"
command_args="director daemon run"
command_user="${ICINGA2_DIRECTOR_USER:-icingadirector}:${ICINGAWEB2_GROUP:-icingaweb2}"
supervisor="supervise-daemon"
