# community/eturnal/eturnal.initd
#!/sbin/openrc-run

name="${RC_SVCNAME}"
command="/usr/sbin/${name}ctl"
command_args="foreground" # Leave backgrounding to start-stop-daemon.
command_background=true
extra_started_commands="reload"
capabilities="^cap_net_bind_service"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}

reload() {
	ebegin "Reloading $name"
	$command reload
	eend $?
}
