# main/busybox/klogd.initd
#!/sbin/openrc-run

name="busybox klogd"
command="/sbin/klogd"
command_args="${KLOGD_OPTS} -n"
command_user="klogd"
pidfile="/run/klogd.pid"
command_background=true
capabilities="^cap_syslog"

depend() {
	need clock hostname localmount logger
	before net
	keyword -vserver -lxc
}
