# main/busybox/ntpd.initd
#!/sbin/openrc-run

name="busybox $RC_SVCNAME"
command="/usr/sbin/$RC_SVCNAME"
command_args="${NTPD_OPTS:--N -p pool.ntp.org} -n"
pidfile="/run/$SVCNAME.pid"
command_background=true

depend() {
	need net
	provide ntp-client
	use dns
}

