# main/busybox/watchdog.initd
#!/sbin/openrc-run

name="busybox watchdog"
command="/sbin/watchdog"
command_args="${WATCHDOG_OPTS} -F ${WATCHDOG_DEV}"
pidfile="/run/watchdog.pid"
command_background=true

depend() {
	need dev
	after hwdrivers
}

start_pre() {
	if ! [ -n "$WATCHDOG_DEV" ]; then
		eerror "WATCHDOG_DEV is not set"
		return 1
	fi
}
