# main/musl-nscd/musl-nscd.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/nscd"
command_args="$NSCD_OPTS"
pidfile="/run/$SVCNAME.pid"
command_background=true

start_pre() {
        mkdir -p /var/run/nscd
}
