# testing/rinetd/rinetd.initd
#!/sbin/openrc-run

description="TCP/UDP port redirector"
cfgfile=${cfgfile:-/etc/rinetd.conf}
command=${command:-/usr/sbin/rinetd}
command_args="-f -c $cfgfile"
required_files="$cfgfile"

depend() {
	need net
	use dns
}

