# community/himitsu-ssh/hissh-agent.initd
#!/sbin/openrc-run

description="himitsu SSH agent"

command="/usr/bin/hissh-agent"
supervisor=supervise-daemon

depend() {
	need himitsud
}
