# testing/lotide/hitide.initd
#!/sbin/openrc-run

name=hitide
description="Lotide UI"

export BACKEND_HOST
export FRONTEND_URL
export HITIDE_PORT

command="/usr/bin/hitide"
command_background=true
command_user="hitide:hitide"
pidfile="/run/${RC_SVCNAME}.pid"
error_log="/var/log/hitide.log"

depend() {
	need localmount net
	after firewall lotide
}

start_pre() {
	checkpath -f -m 0640 -o "$command_user" "$error_log"
}
