# community/hexagonrpcd/hexagonrpcd-adsp-rootpd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Hexagonrpcd ADSP RootPD"
description="Daemon to support Qualcomm Hexagon ADSP virtual filesystem for RootPD"

command=/usr/bin/hexagonrpcd
command_args="-f /dev/fastrpc-adsp -d adsp -R ${hexagonrpcd_fw_dir}"
command_user=fastrpc:fastrpc
