# community/multicast-relay/multicast-relay.initd
#!/sbin/openrc-run

name="multicast-relay"
description="Relay broadcast and multicast packets between interfaces"

command="/usr/bin/multicast-relay.py"
command_args="--foreground --interfaces $INTERFACES $EXTRA_PARAMETERS"
command_background=true

pidfile="/run/${RC_SVCNAME}.pid"

depend() {
        need net
}

