# community/hd-idle/hd-idle.initd
#!/sbin/openrc-run

description="hd-idle"

pidfile="/run/$RC_SVCNAME.pid"
command="/usr/sbin/hd-idle"
command_args="-l ${LOGFILE} -i 0 -c ${DISK_TYPE} -a ${DISK} -i ${IDLE_TIME}"
command_background=yes

depend() {
	need localmount
	after modules
}
