# testing/emby/emby.initd
#!/sbin/openrc-run

name=emby
command=/usr/bin/mono
command_args="$emby_server -programdata $datadir -ffmpeg $ffmpeg -ffprobe $ffprobe"
command_background=true
pidfile=$emby_pid
start_stop_daemon_args="--user $emby_user --group $emby_group"

depend() {
        use logger dns
        need net
        after firewall
}
