# community/cloud-init/cloud-init-hotplugd.initd
#!/sbin/openrc-run

description="cloud-init hotplugd"

pidfile="/run/$RC_SVCNAME.pid"
command="/usr/sbin/cloud-init-hotplugd"
command_background=yes

depend() {
	before cloud-init
}
