# testing/piler/piler.initd
#!/sbin/openrc-run

description="Piler daemon"

command="/usr/bin/$SVCNAME"
command_background=yes

depend() {
	need net
	after firewall mariadb
}

start_pre() {
	checkpath -d -m 750 -o piler:piler /run/piler
}

