# testing/piler/piler-smtp.initd
#!/sbin/openrc-run

description="Piler SMTP server"

command="/usr/bin/$SVCNAME"
command_background=yes

depend() {
	need net
	after firewall
}

