# non-free/vlmcsd/vlmcsd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name=vlmcsd
command="/usr/bin/vlmcsd"
command_user=vlmcsd
command_args="-l syslog -D"
pidfile="/run/vlmcsd.pid"

depend() {
        use logger dns
        need net
        after firewall
}
