# non-free/vlmcsd/vlmcsd.initd
#!/sbin/openrc-run

command=/usr/bin/vlmcsd
pidfile=/run/vlmcsd/vlmcsd.pid
command_args="-l syslog -p ${pidfile}"
name=vlmcsd
start_stop_daemon_args="--user vlmcsd --group vlmcsd"

depend() {
	need networking
}

start_pre() {
	checkpath -d -o vlmcsd:vlmcsd -m 0775 ${pidfile%/*}
}
