# testing/openrc-exporter/openrc-exporter.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=openrc-exporter
description="Prometheus exporter which exposes metrics on services managed by OpenRC"

command=/usr/bin/openrc-exporter
command_args="$openrc_exporter_opts"
command_user="openrc-exporter:openrc-exporter"

depend() {
	need net
	after firewall
}
