# main/openntpd/openntpd.initd
#!/sbin/openrc-run

name="NTP Server"
command=/usr/sbin/ntpd
command_args="-d $NTPD_OPTS"
command_background=yes
pidfile=/run/$RC_SVCNAME.pid
required_files=/etc/ntpd.conf

depend() {
	need net
	provide ntp-client
	use dns logger
}
