# testing/elementary-settings-daemon/elementary-settings-daemon.initd
#!/sbin/openrc-run

description="Portal service (Elementary Settings Daemon)"

command="/usr/libexec/io.elementary.settings-daemon.xdg-desktop-portal"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
	need dbus wireplumber
	before pipewire-pulse
}
