# community/hfd-service/hfd-service.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="HFS service"
description="A dbus activated service that manages human feedback devices sutch as leds and vibrators on mobile devices"

command=/usr/bin/hfd-service

depend() {
	want dbus
}
