# community/hfd-service/hfd-service.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="HFD system service"
description="Manage human feedback devices such as LEDs and vibrators on mobile devices"

command=/usr/libexec/hfd-service

depend() {
	want dbus
}
