# community/syncplay/syncplay-server.initd
#!/sbin/openrc-run

description="Syncplay Server"
command="/usr/bin/${RC_SVCNAME}"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need localmount net
	after firewall
}

# uses the openrc templates for start()/stop()
