# testing/hickory-dns/hickory-dns.initd
#!/sbin/openrc-run

name="Hickory DNS Server"
description="A safe and secure DNS server with DNSSEC support"

command="/usr/bin/hickory-dns"
command_background=yes
command_user=hickory-dns:hickory-dns
capabilities="^cap_net_bind_service"
pidfile="/run/hickory-dns.pid"

depend() {
	need net
}
