# community/mopidy/mopidy.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Mopidy"
description="Extensible music server"

command=/usr/bin/mopidy
command_args="--config /etc/mopidy/mopidy.conf"

command_user=mopidy:audio

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory --owner "$command_user" --mode 0775 \
		/var/log/mopidy
}
