# community/mopidy/mopidy.initd
#!/sbin/openrc-run

name=mopidy
command="/usr/bin/mopidy"
command_args="--config /etc/mopidy/mopidy.conf"
command_user="mopidy"
start_stop_daemon_args="--background --make-pidfile"
pidfile="/run/mopidy.pid"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory --owner mopidy:audio --mode 0775 \
		/var/log/mopidy
}
