# testing/freshrss/freshrss.crontab
*/15 * * * * php85 -f /usr/share/webapps/freshrss/app/actualize_script.php >> /dev/null 2>&1
