# testing/rebuilderd/rebuilderd-sync.cron
#!/bin/sh

. /etc/conf.d/rebuilderd

for profile in $sync_profiles; do
	echo "Syncing $profile"
	rebuildctl pkgs sync-profile "$profile"
done
