# community/alloy/alloy.initd
#!/sbin/openrc-run

name="alloy"
description="OpenTelemetry Collector distribution with programmable pipelines"

: ${command_user:="alloy"}

command="/usr/bin/alloy"
command_args="$command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
