# community/dex-idp/dex-idp.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="OpenID Connect (OIDC) identity and OAuth 2.0 provider with pluggable connectors"
supervisor="supervise-daemon"
command="/usr/bin/dex-idp"
supervise_daemon_args=" serve /etc/dex.yaml"
command_user="dex-idp:dex-idp"

depend() {
	need net
	after firewall
}
