# core/openrc/hostname.initd
#!/sbin/runscript

description="Sets the hostname of the machine."

depend() {
	keyword noprefix
}

start() {
	if [ -f /etc/hostname ] ; then
		opts="-F /etc/hostname"
	else
		opts="localhost"
	fi
	ebegin "Setting hostname"
	hostname $opts
	eend $?
}
