# testing/filebrowser/filebrowser.initd
#!/sbin/openrc-run

description="filebrowser server"
command="/usr/bin/filebrowser"
command_args="$command_args"
command_background=true
command_user="filebrowser:filebrowser"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
    need net
    after firewall
}
