# community/cgmanager/cgmanager.initd
#!/sbin/openrc-run

name="cgroups manager"
description="Daemon to manage cgroups"

: ${retry:="TERM/45"}

command=/usr/sbin/cgmanager
# $OPTIONS is here for backward compatibility only
command_args="${command_args:-$OPTIONS} -m name=system"
command_background=yes
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need cgroups
	before cgproxy
}
