# community/cgmanager/cgproxy.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="cmanager system-bus proxy"
description="D-Bus proxy for cgmanager"

command=/usr/sbin/cgproxy
command_args="$OPTIONS"
retry=TERM/45

depend() {
	need cgmanager dbus
}
