# community/cgmanager/cgproxy.initd
#!/sbin/openrc-run

name="cmanager system-bus proxy"
description="D-Bus proxy for cgmanager"

command=/usr/sbin/cgproxy
# $OPTIONS is here for backward compatibility only
command_args="${command_args:-$OPTIONS}"
retry=TERM/45

depend() {
	need cgmanager dbus
}
