# community/cgmanager/cgproxy.initd
#!/sbin/openrc-run

name="cmanager system-bus proxy"
description="D-Bus proxy for cgmanager"

: ${retry:="TERM/45"}

command=/usr/sbin/cgproxy
# $OPTIONS is here for backward compatibility only
command_args="${command_args:-$OPTIONS}"
command_background=yes
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need cgmanager dbus
}
