# community/incus/incus-agent.initd
#!/sbin/openrc-run

description="Incus Agent"
command="/usr/sbin/incus-agent"
command_args="${INCUS_AGENT_OPTIONS}"
command_background="true"
pidfile="/run/incus/${RC_SVCNAME}.pid"

depend() {
	before cloud-init cloud-init-local
}

start_pre() {
	checkpath --directory "${pidfile%/*}" --mode 0750
}
