# testing/csync2/csync2.initd
#!/sbin/openrc-run

name="$SVCNAME"
pidfile="/var/run/$SVCNAME.pid"
command="/usr/sbin/$SVCNAME"
command_args="-ii"
command_background="yes"
required_files="/etc/csync2/csync2.cfg /etc/csync2/csync2_ssl_key.pem /etc/csync2/csync2_ssl_cert.pem"

depend() {
    need net localmount
    after firewall
}
