# main/bmd-tools/bmd-streamer.initd
#!/sbin/openrc-run

depend() {
	after net firewall chrony
}

description="Blackmagic streamer"
command="/usr/bin/bmd-streamer"
command_background="yes"
pidfile="/var/run/bmd-streamer.pid"
