# community/lightdm/lightdm.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Display Manager"
description="Lightweight $name"

command=/usr/bin/lightdm

depend() {
	need localmount dbus
}

start_pre() {
	checkpath --owner lightdm:lightdm --directory /run/lightdm
}

