# testing/postfwd/postfwd.initd
#!/sbin/openrc-run

command="/usr/sbin/postfwd"
pidfile="/run/postfwd.pid"
command_args="--file ${POSTFWD_CONFIG}
--interface=${POSTFWD_LISTEN}
--port=${POSTFWD_PORT}
--user=${POSTFWD_USER}
--group=${POSTFWD_GROUP}
--pidfile ${pidfile}
${POSTFWD_OPTS}
"
command_args_background="--daemon"

depend() {
	need net
}

