# community/grommunio-gromox/gromox-requeue.cron
#!/bin/sh
# Periodic requeuing of failed messages in delivery queue

queue="/var/lib/gromox/queue"
j=0
for i in "$queue"/save/*; do
    [ -f "$i" ] || continue
    mv "$i" "$queue"/mess/0$j
    j=$(($j+1))
    /usr/bin/logger -t gromox-requeue "requeued $i"
done
