# community/grommunio-gromox/gromox-event.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
description="Gromox event service"
description_reload="Reload configuration without exiting"

command="/usr/libexec/gromox/event"

extra_started_commands="reload"

reload() {
    ebegin "Reloading configuration"
    $supervisor $RC_SVCNAME --signal HUP
    eend $?
}
