# community/grommunio-gromox/gromox-cleaner.cron
#!/bin/sh
# Periodic purge of softdeleted items in Gromox mailboxes

# Load parameters
[ -f /etc/gromox/gromox.cfg ] && . /etc/gromox/gromox.cfg

# Fallback if gromox.cfg has no value
: ${softdelete_purgetime:="30d"}

/usr/sbin/gromox-mbop foreach.here.mb \( purge-softdelete -t ${softdelete_purgetime} -r / \) \( purge-datafiles \)
