# community/grommunio-gromox/gromox-delivery.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name=$RC_SVCNAME
description="Gromox mail spooler"
description_reload="Reload configuration without exiting"

command="/usr/libexec/gromox/delivery"

extra_started_commands="reload"

depend() {
    need gromox-delivery-queue
    after mariadb
}

reload() {
    ebegin "Reloading configuration"
    $supervisor $RC_SVCNAME --signal HUP
    eend $?
}
