# community/vdr/vdr.initd
#!/sbin/openrc-run

depend() {
	use logger
}

command=/usr/bin/vdr
command_args=
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
command_user="vdr:vdr"
