# testing/strongswan/strongswan.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting StrongSwan"
	ipsec start
	eend $?
}

stop () {
	ebegin "Stopping StrongSwan"
	ipsec stop
	eend $?
}
restart() {
	ebegin "Restarting StrongSwan"
	svc_stop
	sleep 2
	svc_start
	eend $?
}

status() {
	ebegin "StrongSwan Status (verbose):"
	ipsec statusall
	eend $?
}
