# main/sntpc/sntpc.cron
#!/bin/sh

# cron job for automatic ntp updates
sleep $(expr $RANDOM % 7200)
exec /etc/init.d/sntpc restart

