# testing/shairport-sync/shairport-sync.initd
#!/sbin/openrc-run

description="Implements a synchronous (multi-room-capable) AirPlay receiver"
command="/usr/bin/shairport-sync"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="
	--stdout /var/log/$RC_SVCNAME.log
	--stderr /var/log/$RC_SVCNAME.log"

depend() {
	need net localmount avahi-daemon
}
