# testing/shairport-sync/shairport-sync.initd
#!/sbin/openrc-run

description="Implements a synchronous (multi-room-capable) AirPlay receiver"
name="shairport-sync $RC_SVCNAME"
command="/usr/bin/shairport-sync"
command_args="-d"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net localmount
	use avahi-daemon
}
