# testing/pd-mapper/pd-mapper.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
name="pd-mapper"
description="Qualcomm PD mapper service"
command="/usr/bin/pd-mapper"

depend() {
	want qrtr-ns
}
