# testing/calibre/calibre.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}
start()  {
	ebegin "Starting Calibre Server"
	start-stop-daemon --start --quiet --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_ARGS
	eend $?
}  
 
stop() {  
	start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --pidfile $PIDFILE  
	rm -f $PIDFILE  
}  
