# community/dotool/dotoold.initd
#!/sbin/openrc-run

description="dotool daemon"

command="/usr/bin/dotoold"
command_background="yes"
command_user="dotool:input"
pidfile="/run/$RC_SVCNAME.pid"
