# testing/snapcast/snapcast-server.initd
#!/sbin/openrc-run

name="$RC_SVCNAME"
command="/usr/bin/snapserver"
command_args="${USER_OPTS:---user snapcast:snapcast} $SNAPSERVER_OPTS"
command_background=yes
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
