# community/go2rtc/go2rtc.initd
#!/sbin/openrc-run

name="go2rtc"
description="Camera streaming application"
command="/usr/bin/go2rtc"
command_args="-config /etc/go2rtc/go2rtc.yaml"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

output_log="/var/log/go2rtc.log"
error_log="/var/log/go2rtc.log"

depend() {
	need net
	after firewall
}
