# testing/darkice/darkice.initd
#!/sbin/runscript

DARKICE_BIN="/usr/bin/darkice"
DARKICE_OPTS="-c /etc/darkice/darkice.cfg"


depend() {
	need net
	after firewall
}

start() { 
	ebegin "Starting DarkIce"
	start-stop-daemon --start --quiet --background --exec  $DARKICE_BIN \
		-- $DARKICE_OPTS > /dev/null
        eend $?
}


stop() {
	ebegin "Stopping DarkIce"
	start-stop-daemon --stop --exec $DARKICE_BIN
	eend $?
}
