# testing/darkice/darkice.initd
#!/sbin/runscript

depend() {
	need net icecast
	after firewall
}

start() { 
	ebegin "Starting DarkIce"
	start-stop-daemon --start --quiet --background \
		--exec ${DARKICE_BIN} -- ${DARKICE_OPTS}
        eend $?
}


stop() {
	ebegin "Stopping DarkIce"
	start-stop-daemon --stop --exec $DARKICE_BIN
	eend $?
}
