# community/qbootctl/qbootctl-openrc.initd
#!/sbin/openrc-run

name="bootctl"
description="Mark current boot as successful in Android A/B interface"

# This marks the current boot as successful, without this after ~7
# reboots A/B devices automatically switch to the inactive slot
# causing a lot of confusion and usually a soft-bricked device
command="qbootctl"
command_args="-m"

start() {
	default_start
}
