# testing/prometheus-snmp-exporter/prometheus-snmp-exporter.initd
#!/sbin/openrc-run
name="prometheus-snmp-exporter"
description="prometheus monitoring system & time series database"
supervisor=supervise-daemon
command=/usr/bin/prometheus-snmp-exporter
command_args="$command_args"
command_user="prometheus:prometheus"

