# community/plocate/plocate.cron
#!/bin/sh
set -e

ulimit -n 4096
updatedb
