# testing/shutdown-clear-machine-id/clear-machine-id.shutdown.initd
#!/sbin/openrc-run

description="Remove dbus machine-id files on shutdown"

start() {
	rm -f /etc/machine-id
	rm -f /var/lib/dbus/machine-id
}
