# testing/haveged/haveged.initd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting haveged"
	start-stop-daemon --start --quiet --exec /usr/sbin/haveged -- ${HAVEGED_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping haveged"
	start-stop-daemon --stop --quiet --pidfile /var/run/haveged.pid
	result=$?
	eend $result
}
