# testing/pully/pully.initd
#!/sbin/openrc-run

name="Pully"
description="Gitlab Webhook receiver"

: ${cfgfile:="/etc/pully.conf"}
: ${command_user:="root:root"}

command="/usr/bin/pully"
command_args="$cfgfile"

required_files="$cfgfile"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f --owner "$command_user" --mode "0644" $logfile
}
