# community/dcc/dcc.initd
#!/sbin/openrc-run

# dcc init.d file for alpine linux.
# edit /var/dcc/dcc_conf instead of passing options
name=dcc
command="/var/dcc/libexec/rcDCC"
command_args="start"
command_background="no"

start_stop_daemon_args="--user dcc:dcc"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}

start_pre() {

}
