# testing/i3bar-river/i3bar-river.initd
#!/sbin/openrc-run

command="/usr/bin/i3bar-river"
supervisor="supervise-daemon"
output_logger="logger -t '${RC_SVCNAME}' -p daemon.info"
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
