# unmaintained/gradm/grsec-rbac.initd
#!/sbin/openrc-run

start() {
	ebegin "Enabling grsecurity RBAC policy"
		gradm -E
	eend $?
}

stop() {
	ebegin "Disabling grsecurity RBAC policy"
		gradm -D
	eend $?
}

