# testing/cloud-init/cloud-init.initd
#!/sbin/runscript

depend() {
  need net
  after firewall
  after cloud-init-local
  before cloud-config
  provide cloud-init
}

start() {
  cloud-init init
  eend 0
}
