# testing/telegraf/telegraf.initd
#!/sbin/openrc-run

# Sample init.d file for alpine linux.

name="telegraf"

config_dir="/etc/telegraf/telegraf.d"
config_file="/etc/telegraf/telegraf.conf"

command="/usr/bin/$name"
command_args="-config-directory ${config_dir} -config ${config_file} ${telegraf_opts}"
command_background="yes"
command_user=telegraf:telegraf

make_pidfile=yes
pidfile="/run/$name.pid"
