# community/pipewire/pipewire-pulse.initd
#!/sbin/openrc-run

description="PulseAudio compatibility layer for PipeWire"

command="/usr/bin/pipewire-pulse"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

export PIPEWIRE_PULSE_NOTIFICATION_FD=3
ready=fd:3

depend() {
	need pipewire
	want wireplumber
}
