# community/pipewire/pipewire.initd
#!/sbin/openrc-run

description="Multimedia processing graphs"

command="/usr/bin/pipewire"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

export PIPEWIRE_NOTIFICATION_FD=3
ready=fd:3

depend() {
	need dbus
	want wireplumber
	before pipewire-pulse
}
