# testing/usb-signaller/usb-signaller.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name=usb-signaller
description="usb-signaller USB gadget controller"
command="/usr/bin/usb-signaller"
output_logger="logger -t usb-signaller -p daemon.info"
error_logger="logger -t usb-signaller -p daemon.warning"

depend() {
	need dbus
}
