# testing/oo7/oo7-portal.initd
#!/sbin/openrc-run

name="OO7 Portal"
description="An implementation of org.freedesktop.impl.portal.Secret"

supervisor=supervise-daemon
command="/usr/libexec/oo7-portal"

depend() {
	need dbus
}
