# main/openvswitch/ovs-modules.initd
#!/sbin/openrc-run

MODULES="openvswitch"

depend() {
	before ovsdb-server
}

start() {
	ebegin "Loading Openvswitch kernel modules"
	for mod in $MODULES; do
		modprobe -q $mod
	done
	eend $?
}

stop() {
	# modules cannot be removed without a reboot
	return 0
}
