# community/mmsd-tng/mmsd-tng.initd
#!/sbin/openrc-run

description="mmsd-tng daemon"

command="/usr/bin/mmsdtng"
supervisor=supervise-daemon
