# community/himitsu-secret-service/hisecrets-agent.initd
#!/sbin/openrc-run

description="himitsu secrets agent daemon"

command="/usr/bin/hisecrets-agent"
supervisor=supervise-daemon

depend() {
	need himitsud
}
