# testing/anubis/anubis.initd
#!/sbin/openrc-run

name=anubis
description="Anubis HTTP defense proxy"
supervisor="supervise-daemon"
command="/usr/bin/anubis"

depend() {
	need net
	after firewall
}
