# community/lua-turbo/lua-turbo.initd
#!/sbin/openrc-run

name="Turbo Lua"
command="/usr/bin/luajit"
command_args="${TURBO_SCRIPT}"
pidfile="/var/run/turbo.pid"
start_stop_daemon_args="-u ${TURBO_USER} -g ${TURBO_GROUP}"
command_background="yes"

depend() {
        need net
}
