# testing/qrtr/qrtr-ns.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="QRTR"
description="Qualcomm IPC router name service"

# -f should be in command_args_foreground, but then the order is wrong (1 -f)
command="/usr/bin/qrtr-ns"
command_args="-f 1"
