# testing/thelounge/thelounge.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="thelounge"
command="/usr/bin/thelounge"
command_args="start"
command_user="thelounge:thelounge"
command_background=true
pidfile="/run/thelounge.pid"

depend() {
	need net localmount
	after firewall
}
