# community/tuned/tuned-ppd.initd
#!/sbin/openrc-run

description="power-profiles-daemon compatibility for tuned"
supervisor=supervise-daemon
command=/usr/sbin/tuned-ppd
command_args_foreground="-l"

depend() {
	need dbus
	need tuned
}

