# community/yarr/yarr.initd
#!/sbin/openrc-run

name=yarr
command="/usr/bin/yarr"
command_user="yarr:yarr"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net localmount
	after firewall
}

start_pre() {
    checkpath -f "$output_log" -o "$command_user"
    checkpath -f "$error_log" -o "$command_user"
}
