# community/htpdate/htpdate.initd
#!/sbin/openrc-run

name="HTTP Time Protocol service"
command=/usr/sbin/htpdate
pidfile=/run/htpdate.pid
command_args="-D -i $pidfile $HTP_OPTIONS $HTP_SERVERS"

depend() {
	need net
	provide ntp-client
	use dns logger
}
