# testing/musl-nscd/musl-nscd.initd
#!/sbin/openrc-run

name="$SVCNAME"
pidfile="/var/run/$SVCNAME.pid"
command="/usr/sbin/$SVCNAME"
command_args="${NSCD_OPTS:-d -p ${pidfile}}"
