# testing/fenrir/fenrir.initd
#!/sbin/openrc-run

name="fenrir"
description="Fenrir screen reader"
command=/usr/bin/fenrir
pidfile="/run/$RC_SVCNAME.pid"

extra_started_commands="reload"
description_reload="Reload fenrir configuration"

depend() {
	need udev-settle
}

reload() {
	ebegin "Reloading fenrir configuration"
	start-stop-daemon --signal HUP --pidfile "$pidfile"
	eend $?
}
