# community/ifstate/ifstate.initd
#!/sbin/openrc-run

description="IfState"

depend() {
	need localmount
	after bootmisc hwdrivers modules
	provide net
	keyword -shutdown -stop -jail -prefix -vserver -docker
}

start() {
	ebegin "Starting IfState"
	ifstatecli -q apply
}

restart() {
	ebegin "Restarting IfState"
	ifstatecli -q apply
}

status() {
	ifstatecli check
}
