# testing/cronie/cronie.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/crond"
command_args="$CRON_OPTS"

depend() {
	provide cron
}
