# testing/logcheck-busybox/logcheck.initd
#!/sbin/openrc-run

name="$RC_SVCNAME"
description="Provide @reboot functionality for logcheck in Busybox"
command="/etc/periodic/hourly/logcheck reboot"

depend() {
	need localmount
}

start() {
        einfo "Logcheck: cleaning all dead.letter"
        rm -f /var/lib/logcheck/dead.letter
        rm -f /srv/lxc/*/rootfs/var/lib/logcheck/dead.letter
}

stop() {
        return 0
}
