# main/ppp/pppd.initd
#!/sbin/openrc-run

name="pppd"
command="/usr/sbin/${SVCNAME}"
command_args="nodetach ${PPPD_OPTS}"
pidfile="/var/run/$SVCNAME.pid"
command_background=true

start_pre() {
	checkpath --directory /var/run/ppp
	modprobe pptp
}

