# testing/rtkit/rtkit.initd
#!/sbin/openrc-run

depend() {
	need dbus
}

start() {
	ebegin "Starting rtkit"
	start-stop-daemon --start --quiet --exec /usr/libexec/rtkit-daemon \
	--background --pidfile /var/run/rtkit.pid --make-pidfile \
	-- ${RTKIT_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping rtkit"
	start-stop-daemon --stop --quiet --pidfile /var/run/rtkit.pid --exec /usr/libexec/rtkit-daemon
	eend $?
}
