# testing/tlp/tlp.initd
#!/sbin/openrc-run

description="Initialize tlp"
extra_started_commands="force_reload"

depend() {
	need localmount
}

status() { tlp-stat -s; }
# TODO: keep the following in one line
start() { tlp init start; }
stop() { tlp init stop; }
restart() { tlp init restart; }
force_reload() { tlp init force-reload; }
