# testing/tlp/tlp.initd
#!/sbin/openrc-run

description="Initialize tlp"
command=/usr/sbin/tlp

extra_commands="force-reload"

depend() {
	need localmount
}

status() {
	tlp-stat -s
}

# TODO: I don't know how to write this initd
case "$1" in
    start|\
    stop|\
    restart|\
    force-reload)
        tlp init $1
        ;;
esac
