# testing/lizardfs/lizardfs-metalogger.initd
#!/sbin/openrc-run

command=/usr/sbin/mfsmetalogger

depend() {
	need syslog
	after net
}
