# testing/lizardfs/lizardfs-chunkserver.initd
#!/sbin/openrc-run

command=/usr/sbin/mfschunkserver

depend() {
	need syslog
	after net
}

