# main/multipath-tools/multipathd.initd
#!/sbin/openrc-run

description="Multipath monitoring daemon"

command="/usr/sbin/multipathd"
command_args="-d $command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
output_log="/var/log/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME.log"
