# main/multipath-tools/multipath.initd
#!/sbin/openrc-run

description="Initialize multipath"

depend() {
	before checkfs fsck multipathd iscsid lvm
	after modules device-mapper
}

start() {
	ebegin "Loading multipath modules"
	modprobe -aq dm_mod dm-multipath
	eend $?

	ebegin "Activating multipath devices"
	multipath -v0 >/dev/null
	eend $?
}

stop() {
	ebegin "Deactivating multipath devices"
	multipath -v0 -F >/dev/null
	eend $?
}
