# main/multipath-tools/multipathd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/sbin/multipathd
command_args_foreground="-d"
output_log=/var/log/${RC_SVCNAME}.log
error_log=/var/log/${RC_SVCNAME}.log
