# community/hkdm/hkdm.openrc
#!/sbin/openrc-run

command="/usr/sbin/hkdm"

name="HKDM"
description="Hotkey Daemon (For) Mobile"
