# community/dnscrypt-proxy/dnscrypt-proxy.initd
#!/sbin/openrc-run

command="/usr/sbin/dnscrypt-proxy"
command_args="-config /etc/dnscrypt-proxy/dnscrypt-proxy.toml -syslog"
command_user="${DNSCRYPT_USER:-dnscrypt}:${DNSCRYPT_GROUP:-dnscrypt}"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"

depend() {
    use net logger
    provide dns
}

start_pre() {
    checkpath -q -d -m 0775 -o "${command_user}" /var/log/"${RC_SVCNAME}"
}
