# testing/dkimproxy/dkimproxy.in.initd
#!/sbin/runscript

command="/usr/sbin/${RC_SVCNAME}"
command_args="${OPTS}"
pidfile="/var/run/${RC_SVCNAME}.pid"
description="SMTP proxy that verifies signature on incomming messages"

depend() {
	need net
	after firewall
}

