# testing/woodpecker/woodpecker-agent.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Woodpecker CI agent"

command="/usr/bin/woodpecker-agent"
command_args="-env-file ${env_file:-/etc/woodpecker/agent.conf}"
command_user="woodpecker:woodpecker"
directory="/var/lib/woodpecker"

error_log="/var/log/woodpecker/agent.log"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -m 755 -o "woodpecker:woodpecker" /var/lib/woodpecker
	checkpath -f -m 644 -o "woodpecker:woodpecker" /var/log/woodpecker/agent.log
}
