# testing/tofutf/tofutfd.initd
#!/sbin/openrc-run

name="tofutf server"

: ${command_user:="tofutf:tofutf"}

command="/usr/bin/tofutfd"
command_args="--secret $(cat $secret_file) $command_args"
command_background="yes"
error_log="/var/log/tofutfd.log"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

start_pre() {
        checkpath -f --owner "$command_user" --mode 0644 \
                $error_log
}
