# community/celery/celery.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

description="celery queue worker"

: ${CELERY_USER:="celery"}
: ${CELERY_GROUP:="$(id -gn $CELERY_USER)"}

pidfile="/run/$RC_SVCNAME.sd.pid"
supervise_daemon_args="-u $CELERY_USER -g $CELERY_GROUP"
command=/usr/bin/celery
command_args="${CELERY_OPTS}"

depend() {
	use net
}
