# testing/exim/exim.initd
#!/sbin/openrc-run

name="$SVCNAME"
pidfile="/var/run/$SVCNAME.pid"
command="/usr/sbin/$SVCNAME"
command_args="${EXIM_OPTS:--bd -q15m}"
required_files="/etc/exim/${SVCNAME}.conf"
extra_started_commands="reload"

depend() {
	need logger
	use antivirus net
	provide mta
}

reload() {
	ebegin "Reloading ${SVCNAME}"
	start-stop-daemon --signal HUP --pidfile ${pidfile} --name ${SVCNAME}
	eend $?
}
