# testing/exim/exim.initd
#!/sbin/openrc-run

pidfile="/run/exim.pid"
command="/usr/sbin/exim"
command_args="${EXIM_OPTS:--bd -q15m}"
required_files="/etc/exim/exim.conf"
extra_started_commands="reload"

depend() {
	need logger
	use antivirus net
	provide mta
}

reload() {
	ebegin "Reloading ${SVCNAME}"
	start-stop-daemon --signal HUP --pidfile ${pidfile} --name ${SVCNAME}
	eend $?
}
