# testing/ddnsc/ddnsc.initd
#!/sbin/openrc-run
name="$SVCNAME"
command="/usr/bin/python3 /usr/share/ddnsc/ddnsc.py"
pidfile="/var/run/$SVCNAME.pid"
command_background=true

depend() {
        need net
        after firewall
}
