# testing/nldev/nldev.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/sbin/$SVCNAME"
command_args="${NLDEV_OPTS:--k}"
pidfile="/var/run/$SVCNAME.pid"
command_background="yes"

depend() {
	provide dev
	need sysfs dev-mount
	before checkfs fsck
	keyword -lxc -vserver
}
