# community/wlsunset/wlsunset.initd
#!/sbin/openrc-run

description="Gamma adjustments for Wayland compositors"

command="/usr/bin/wlsunset"
command_args="-l ${LATITUDE} -L ${LONGITUDE} ${command_args:-}"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
