# community/yambar/yambar.initd
#!/sbin/openrc-run

command="/usr/bin/yambar"
supervisor=supervise-daemon
command_args="--print-pid=3 $command_args"
notify=fd:3

depend() {
	use dbus
}
