# testing/heplify/heplify.initd
#!/sbin/openrc-run

description="Multi-Protocol Go HEP Capture Agent"
command=/usr/bin/heplify
command_args=${CMD_OPTS:="-i lo"}
command_background=yes
pidfile="/run/homer/${RC_SVCNAME}.pid"

depend() {
	need net
	after firewall
}

