# testing/lightdm/lightdm.initd
#!/sbin/runscript

depend() {
	need localmount dbus
}

description="Lightweigh display manager"
command="/usr/sbin/lightdm"
command_background="yes"
pidfile="/var/run/lightdm.pid"

start_pre() {
	checkpath --owner lighdm:lightdm --directory /var/run/lightdm
}

