# testing/lightdm/lightdm.initd
#!/sbin/runscript

depends() {
	needs localmount dbus
}

description="Lightweigh display manager"
command="/usr/sbin/lightdm"
command_background="yes"
pidfile="/var/run/lightdm.pid"

