# testing/lightdm/lightdm.initd
#!/sbin/openrc-run

depend() {
	need localmount dbus
}

description="Lightweight Display Manager"
command="/usr/bin/lightdm"
command_background="yes"
pidfile="/run/lightdm.pid"

start_pre() {
	checkpath --owner lightdm:lightdm --directory /run/lightdm
}

