# testing/grafana/grafana-server.initd
#!/sbin/openrc-run

pidfile=/run/grafana.pid
start_stop_daemon_args="--user grafana --chdir /usr/share/grafana"
command=/usr/bin/grafana-server
command_args="${EXTRA_OPTS:--config /etc/grafana/grafana.ini}"
command_background="yes"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -o grafana:grafana -m755 /var/lib/grafana
}
