# community/aws-session-manager-plugin/ssmcli.initd
#!/sbin/openrc-run
description="Initialize Amazon SSMCLI"
command="/usr/bin/ssmcli"

start() {
	ebegin "Starting SSMCLI"
	start-stop-daemon --start \
		--exec ${command} ${command_args} \
		-b --make-pidfile \
		--pidfile "${pidfile}"
	eend $?
}

stop() {
	ebegin "Stopping SSMCLI"
	start-stop-daemon --stop \
		--exec ${command} \
		--pidfile "${pidfile}"
	eend $?
}
