# testing/ergo/ergo.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="ergo"
command="/usr/bin/ergo"
command_args="run --conf ${ERGO_CONFIGFILE:-'/etc/ergo/ircd.yaml'}"
command_user="ergo:ergo"
command_background=true
pidfile="/run/ergo.pid"

depend() {
	need net localmount
	after firewall
}
