# testing/sssd/sssd.initd
#!/sbin/openrc-run

command="/usr/sbin/sssd"
command_args="-D $SSSD_OPTIONS"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend(){
	need localmount netmount clock
}
