# testing/llmnrd/llmnrd.initd
#!/sbin/openrc-run

name=llmnrd
command="/usr/sbin/llmnrd"
command_args="${LLMNRD_OPTS}"
command_background="yes"

pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
