# testing/q6voiced/q6voiced.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="q6voiced"
description="Enable q6voice audio when call is performed with oFono"

# Note: q6voice_card/q6voice_device need to be set in /etc/conf.d/q6voiced
command="/usr/bin/q6voiced"
command_args="hw:${q6voice_card},${q6voice_device}"
supervise_daemon_args="--user nobody --group audio"

depend() {
	need dbus
}
