# main/bind/bind.initd
#!/sbin/runscript

NAME=named
DAEMON=/usr/sbin/$NAME

depend() {
	need net
	after firewall
	use logger
	provide dns
}

start() {
	ebegin "Starting ${NAME}"
	start-stop-daemon --start --quiet --background \
		--exec ${DAEMON} --nicelevel ${NICELEVEL} \
		-- -u ${USER} -n ${CPU} ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"
	start-stop-daemon --stop --quiet --pidfile $PIDFILE
	eend $?
}
