# community/mpd/mpd.user.initd
#!/sbin/openrc-run

name="Music Player Daemon"
description="A daemon for playing music"

command=/usr/bin/mpd
command_args_foreground="--no-daemon"
