# community/usb-moded/usb-moded.initd
#!/sbin/openrc-run

name=usb-moded
description="usb-moded USB gadget controller"
command="/usr/bin/usb_moded"
command_background=yes

depend() {
	need dbus
}
