# community/seatd/seatd.initd
#!/sbin/openrc-run

name="seatd"
description="Seat management daemon"

command="/usr/bin/seatd"
command_args="-g seat -l ${loglevel:-error} ${command_args:-}"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"

if [ "$loglevel" != 'silent' ]; then
	: ${error_logger="logger -t $name -p daemon.${loglevel:-error} >/dev/null 2>&1"}
fi
