# community/seatd/seatd.initd
#!/sbin/openrc-run

name="seatd"
description="Seat management daemon"

command="/usr/bin/seatd"
command_args="-g seat ${command_args:-}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

error_log="/var/log/seatd.log"
