# community/seatd/seatd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=$RC_SVCNAME
description="Seat management daemon"

command="seatd"
command_args="-g seat"

error_log="/var/log/seatd.log"
