# unmaintained/arpalert/arpalert.initd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting arpalert"
	start-stop-daemon --quiet --start --pidfile /var/run/arpalert.lock --exec /usr/sbin/arpalert -- -d
	eend $?
}

stop() {
	ebegin "Stopping arpalert"
	start-stop-daemon --quiet --stop --pidfile /var/run/arpalert.lock
	eend $?
}

