# main/aconf/aconf.initd
#!/sbin/openrc-run

# init.d file for aconf
# Copyright (c) 2013-2014 Kaarle Ritvanen

name=aconf
command=/usr/sbin/aconfd
pidfile=/var/run/$name.pid

depend() {
	need net
	after firewall
}

stop() {
	ebegin "Stopping $name"
	start-stop-daemon --stop --pidfile $pidfile --signal INT
	eend $?
}
