# community/gitlab-runner/gitlab-runner.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

description="GitLab CI Runner"
name="Gitlab Runner"

command="/usr/bin/gitlab-runner"
command_args="run ${GITLAB_RUNNER_OPTS}"
command_user="${GITLAB_RUNNER_USER}:${GITLAB_RUNNER_GROUP}"

depend() {
	need net
}
