# testing/pnp4nagios/npcd
#!/sbin/runscript

depend() {
        use net
        need dbus
        before nfs
        after logger
	}

start() {
        ebegin "Starting npcd"
        start-stop-daemon --start --quiet --exec /usr/bin/npcd -- -d -f /etc/pnp4nagios/npcd.cfg
        eend $?
	}

stop() {
       ebegin "Stopping npcd"
       start-stop-daemon --stop --quiet --exec /usr/bin/npcd
       eend $?
	}

