# testing/sleep-inhibitor/sleep-inhibitor.initd
#!/sbin/openrc-run
supervisor="supervise-daemon"
description="Starting Sleep Inhibitor"
command="/usr/bin/sleep-inhibitor"

depend() {
        after bootmisc
}
