# community/tvheadend/tvheadend.initd
#!/sbin/openrc-run

command="/usr/bin/tvheadend"
command_args="-f -c $TVHCONF -u tvheadend -g video -p $TVHPID $TVHARGS"
retry="TERM/30/KILL/5"
pidfile=$TVHPID
name="tvheadend-git"
required_dirs=$TVHCONF

depend() {
	use logger dns
	need net
	after firewall
}
