# testing/rethinkdb/rethinkdb.initd
#!/sbin/openrc-run

name=rethinkdb
command="/usr/bin/rethinkdb"
command_args="--config-file /etc/rethinkdb/default.conf"

depend() {
        use logger dns
        need net
        after firewall
}

