# testing/interception-tools/udevmon.initd
#!/sbin/openrc-run

command=/usr/bin/udevmon
command_args="-c /etc/interception/udevmon.yaml"
pidfile=/run/udevmon.pid
command_background=true

depend() {
	want dev-settle
	after dev-settle
}
