# community/phodav/spice-webdavd.initd
#!/sbin/openrc-run

description="Spice guest utility for file sharing via WebDAV"

command="/usr/sbin/spice-webdavd"
command_args="--port 9843"

pidfile="/run/${SVCNAME}.pid"
command_background="yes"

depend() {
	need dbus 
	use avahi
}
