# testing/speakersafetyd/speakersafetyd.initd
#!/sbin/openrc-run

pidfile="/run/speakersafetyd.pid"
command="/usr/bin/speakersafetyd"
command_args="-c /usr/share/speakersafetyd/ -b /var/lib/speakersafetyd/blackbox -m 7"
supervisor="supervise-daemon"
	
