# testing/amavisd-new/amavisd.init
#!/sbin/runscript

depend() {
	need net
	use logger
	use antivirus
	before mta
}

start() {
	ebegin "Starting amavisd-new"
	start-stop-daemon --start --quiet --name amavisd --pidfile /var/amavis/amavisd.pid \
	 --exec /usr/sbin/amavisd
	eend $?
}

stop() {
	ebegin "Stopping amavisd-new"
	start-stop-daemon --stop --quiet --pidfile /var/run/amavis/amavisd.pid
	eend $?
}
