# testing/bird-lg-go/bird-lgproxy-go.initd
#!/sbin/openrc-run

name="${RC_SVCNAME}"
pidfile="/run/$name.pid"

command="/usr/bin/bird-lgproxy-go"
: ${command_user:="bird-lgproxy-go:bird"}
command_background=true

depend() {
	need bird
}
