# testing/bird-lg-go/bird-lg-go.initd
#!/sbin/openrc-run

name="${RC_SVCNAME}"
pidfile="/run/$name.pid"

command="/usr/bin/bird-lg-go"
: ${command_user:="bird-lg-go:bird-lg-go"}
command_background=true

depend() {
	need net
	use bird bird-lg-proxy
	after firewall
}
