# main/bluez/bluetooth.initd
#!/sbin/runscript

name="Bluetooth"
command="/usr/lib/bluetooth/bluetoothd"
pidfile="/var/run/bluetoothd.pid"
start_stop_daemon_args="--background --make-pidfile"


depend() {
	after coldplug
	need dbus localmount
}

start_post() {
	local adapter
	eindent
	for adapter in $AUTO_ENABLE; do
		ebegin "Enabling $adapter"
		hciconfig $adapter up
		eend
	done
	eoutdent
	return 0
}

