# main/bluez/bluetooth.initd
#!/sbin/openrc-run

name="Bluetooth"
command="/usr/lib/bluetooth/bluetoothd"
pidfile="/var/run/bluetoothd.pid"
start_stop_daemon_args="--background --make-pidfile"


depend() {
	after coldplug
	need dbus localmount hostname
}

