# main/bluez/bluetooth.initd
#!/sbin/runscript

name="Bluetooth"
command="/usr/sbin/bluetoothd"

depend() {
	after coldplug
	need dbus localmount
}

