# community/tpm2-abrmd/tpm2-abrmd.initd
#!/sbin/openrc-run

command="/usr/sbin/tpm2-abrmd"
command_args="${TPM2_ABRMD_ARGS} --logger=syslog"
command_background=1
command_user="tss:tss"
description="TPM2 Access Broker & Resource Management Daemon implementing the TCG spec"
pidfile="/run/tpm2-abrmd.pid"

depend() {
	use logger
}
