# testing/thanos/thanos.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=thanos
description="Highly available Prometheus setup with long term storage capabilities"
command="/usr/bin/thnaos"
command_args="$THANOS_EXEC $THANOS_OPTS"

rc_ulimit="${THANOS_ULIMIT:--n 65536}"

start_stop_daemon_args="--user prometheus:prometheus"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall prometheus
}
