# testing/captagent/captagent.initd
#!/sbin/openrc-run

command=/usr/bin/captagent
pidfile=/var/run/${SVCNAME}.pid
conf=/etc/captagent/${SVCNAME}.xml

#start_stop_daemon_args="--make-pidfile"

depend() {
	need net
	after firewall
}

start_pre() {
	#checkpath --directory ${pidfile%/*}
	command_args="-d -f $conf"
}

