# testing/captagent/captagent.initd
#!/sbin/openrc-run

# init.d file for captagent
# Copyright (c) 2017 Kaarle Ritvanen

name=captagent
command=/usr/bin/$name
command_args=-d
pidfile=/var/run/$name.pid

depend() {
	need net
}
