# testing/lcdproc/lcdproc.initd
#!/sbin/openrc-run

command="/usr/bin/$SVCNAME"
command_args="${LCDD_OPTS:--f}"
command_background=1
pidfile="/run/$RC_SVCNAME.pid"

depend() {
        need net
        use LCDd
}
