# main/mdevd/mdevd-init.initd
#!/sbin/openrc-run
# Copyright 2022 Laurent Bercot for Alpine Linux
# Distributed under the terms of the ISC License.

description="coldplug trigger for mdevd"

depend() {
	need mdevd
	before checkfs fsck
	keyword -containers -vserver -lxc
}

start() {
	ebegin "Scanning hardware for mdevd"
	mdevd-coldplug
	eend $?
}
