# community/openbgpd/openbgpd.initd
#!/sbin/openrc-run
name=$RC_SVCNAME
description="Free implementation of the Border Gateway Protocol, Version 4"
command=/usr/sbin/bgpd
command_args="${BGPD_OPTS}"

extra_started_commands="reload"

depend() {
	use net
	use logger
}

reload() {
	${command} -n || return 1
	ebegin "Reloading bgpd"
	/usr/sbin/bgpctl reload
	eend $?
}
