# testing/ovn/ovn-northd.initd
#!/sbin/openrc-run

command="/usr/share/ovn/scripts/ovn-ctl"
command_args="start_northd --ovn-manage-ovsdb=no --no-monitor $OVN_CTL_OPTS"
command_background="yes"
pidfile="/run/ovn/$RC_SVCNAME.pid"
description="Open Virtual Network central control daemon"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory "${pidfile%/*}"
}

stop() {
	$command stop_northd --no-monitor
}
