# community/oo7/oo7-server.initd
#!/sbin/openrc-run

name="OO7 server"
description="D-Bus Secret Service provider"

supervisor=supervise-daemon
command="/usr/libexec/oo7-daemon"

depend() {
	need dbus
}
