# testing/supervisor/supervisord.initrd
#!/sbin/openrc-run

name="Supervisor"
command="/usr/bin/supervisord"
pidfile="/var/run/supervisord.pid"
configfile="/etc/supervisord.conf"
command_args="--nodaemon --pidfile ${pidfile} --configuration ${configfile}"
command_background="yes"

required_files="${configfile}"
