# testing/hiawatha/hiawatha.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Webserver"
description="Secure and Advanced Hiawatha Webserver"

command=/usr/sbin/hiawatha
command_args_foreground="-d"

depend() {
	after firewall
}
