# testing/hiawatha/hiawatha.initd
#!/sbin/openrc-run

pidfile=/var/run/hiawatha.pid
command=/usr/sbin/hiawatha

depend() {
	after firewall
}
