# testing/hiawatha/hiawatha.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Hiawatha webserver"
description="Secure and advanced webserver"

command=/usr/sbin/hiawatha
command_args_foreground="-d 1>/dev/null"
pidfile=/var/run/hiawatha.pid
healthcheck_timer=60

depend() {
	after firewall
}

stop() {
	local supervise_pid=$(pgrep -f "$supervisor $RC_SVCNAME")

	einfo "Stopping $name"
	kill -TERM $supervise_pid
}
