# community/greenbone-security-assistant/gsad.initd
#!/sbin/openrc-run
#Copyright 2016 Alpine Linux
# Distributed under the terms of the GNU General Public License v2
# $Id$

command_args="${GSAD_LISTEN} ${MANAGER_LISTEN} ${MANAGER_PORT} ${GSAD_EXTRA_ARGS}"
name="Greenbone Security Assistant (GSA)"
command=/usr/bin/gsad

pidfile="/run/gsad.pid"

depend() {
	after bootmisc
	need net
}

