# community/gitea/gitea.initd
#!/sbin/openrc-run

name=gitea
conffile="$GITEA_CUSTOM/conf/app.ini"
command="/usr/bin/gitea"
command_args="web -c $conffile"
start_stop_daemon_args="${GITEA_USER:+--user} $GITEA_USER --env GITEA_CUSTOM=$GITEA_CUSTOM"
pidfile="/var/run/gitea.pid"
command_background="yes"

depend() {
        use logger dns
        need net
        after firewall
}
