# community/gitea/gitea.initd
#!/sbin/openrc-run

name=gitea
command="/usr/bin/gitea"
command_args="web"
start_stop_daemon_args="--user ${GITEA_USER:-gitea} --env GITEA_CUSTOM=${GITEA_CUSTOM:-/var/lib/gitea} --chdir ${GITEA_CUSTOM:-/var/lib/gitea}"
command_background=yes
pidfile="/run/gitea.pid"

depend() {
        use logger dns
        need net
        after firewall
}
