# testing/crab-hole/crab-hole.initd
#!/sbin/openrc-run

name="Crab-Hole DNS Server"
description="A DNS level ad and spy blocker"

command="/usr/bin/crab-hole"
command_background=yes
command_user=crab-hole:crab-hole
capabilities="^cap_net_bind_service"
pidfile="/run/crab-hole.pid"

depend() {
	need net
}

start_pre() {
	checkpath --directory --owner "$command_user" /var/lib/crab-hole
}
