# testing/soapy-sdr-remote/soapy-sdr-remote.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="soapy-sdr-remote"
description="SoapyRemote network server"

command="/usr/bin/SoapySDRServer"
command_args="--bind"

output_log="/var/log/soapy-sdr-remote.log"

depend() {
	need avahi-daemon net localmount
	after firewall
}
