# main/xdm/xdm.initd
#!/sbin/openrc-run

depends() {
	needs localmount
}

description="X Display Manager"
command="/usr/bin/xdm"
command_args="-resources /etc/X11/xdm/Xresources.alpinelinux"
pidfile="/var/run/xdm.pid"

# uses the openrc templates for start()/stop()

