# community/waked/waked.initd
#!/sbin/openrc-run
name="wake daemon"
description="wakes up system at times set by other applications"

command=/usr/bin/waked
command_background=yes
pidfile="/run/$RC_SVCNAME.pid"
