# community/msmtp/msmtp.initd
#!/sbin/openrc-run

: ${MSMTP_INTERFACE:="127.0.0.1"}
: ${MSMTP_PORT:="25"}
: ${MSMTP_COMMAND:="/usr/bin/msmtp -f %F"}

pidfile="/var/run/${SVCNAME}"
command="/usr/bin/msmtpd"
command_args="--interface=$MSMTP_INTERFACE --port=$MSMTP_PORT --command='$MSMTP_COMMAND' "
command_background="yes"
