# testing/prometheus-node-exporter/prometheus-node-exporter.initd
#!/sbin/openrc-run

name="prometheus node exporter"
command="/usr/bin/node_exporter"
command_background="yes"

pidfile="/run/$SVCNAME.pid"
start_stop_daemon_args="--user $prometheus_user:$prometheus_group -mp $pidfile"

depend() {
	need net
	after firewall
}
