# testing/radicale/radicale.initd
#!/sbin/openrc-run

pidfile=/var/run/radicale.pid
start_stop_daemon_args="--user radicale"
command=/usr/bin/radicale
command_args="--foreground"
command_background=yes

depend() {
	use net
	need localmount
}

#start_pre() {
#	checkpath --directory "${pidfile%/*}"
#}
