# testing/dmarc-metrics-exporter/dmarc-metrics-exporter.initd
#!/sbin/openrc-run

name="dmarc-metrics-exporter"

: ${command_user:="prometheus"}

command="/usr/bin/dmarc-metrics-exporter"
command_args="$command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
