# testing/netbird/netbird.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/netbird
command_user=root:netbird

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -m0775 -o root:netbird /var/log/netbird
}
