# testing/iscsi-scst-svn/iscsi-scst.initd
#!/sbin/runscript
#
# Start the iSCSI-SCST Target.
#

depend() {
	need net
	afer firewall
}

start() {
	modprobe -q crc32c
	modprobe -q iscsi-scst
	ebegin "Starting iscsi-scstd"
	start-stop-daemon --start --exec /usr/sbin/iscsi-scstd
	eend 0
}

stop() {
	ebegin "Stopping iscsi-scstd"
	start-stop-daemon --stop --exec /usr/sbin/iscsi-scstd
	eend 0
}
