# community/cronie/cronie.initd
#!/sbin/openrc-run

name="$RC_SVCNAME"
command="/usr/sbin/crond"
command_args="$CRON_OPTS -n"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	provide cron
}
