# testing/upmpdcli/upmpdcli.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="upmpdcli UPnP Media Renderer"
description="UPnP Media Renderer front-end for MPD, the Music Player Daemon"

command=/usr/bin/upmpdcli
command_args="-c /etc/upmpdcli.conf"
command_user=upmpdcli:upmpdcli

depend() {
	need net localmount
	after firewall
}
