# community/vault/vault.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Vault server"
description="Vault is a tool for securely accessing secrets"

VAULT_LOG_FILE="/var/log/${RC_SVCNAME}.log"

command=/usr/sbin/${RC_SVCNAME}
command_args="${vault_opts}"
command_user="${RC_SVCNAME}:${RC_SVCNAME}"

start_pre() {
	checkpath -f -m 0644 -o "$command_user" "$VAULT_LOG_FILE"
}

depend() {
	need net
	after firewall
}
