# main/iperf/iperf.initd
#!/sbin/openrc-run

command="/usr/bin/iperf"
pidfile="/var/run/$SVCNAME/$SVCNAME.pid"

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting $SVCNAME"
	start-stop-daemon --start --make-pid --user ${exec_user:-nobody} \
		--pidfile ${pidfile} --background --exec $command \
		-- --server ${command_args}
	eend
}

