# testing/lxcfs/lxcfs.initd
#!/sbin/openrc-run
command=/usr/bin/lxcfs
command_args="-f -o allow_other /var/lib/lxcfs"
pidfile="/run/lxcfs.pid"
command_background="yes"
retry="TERM/45"
