# main/rrdcollect/rrdcollect.initd
#!/sbin/openrc-run

description="rrdcollect daemon"
pidfile="/var/run/${RC_SVCNAME}.pid"
command="/usr/sbin/rrdcollect"

