# testing/sanlock/sanlock.initd
#!/sbin/openrc-run

name=sanlock
command="/usr/sbin/$name"
command_args="daemon $sanlock_opts"

pidfile="/run/$name/$name.pid"

depend() {
	want wdmd
}
