# testing/entropybroker/entropybroker.initd
#!/sbin/openrc-run

command="/usr/bin/entropy_broker"
command_args="$EB_OPTS"

depend() {
	need net
}
