# community/mumble/murmur.initd
#!/sbin/openrc-run

depend() {
	need localmount net
	after firewall
}

description="Mumble VoIP Server"
command="/usr/bin/murmurd"
command_args="-ini ${INI_CONFIG:-/etc/murmur.ini}"
pidfile="/var/run/murmur/murmur.pid"

start_pre() {
	checkpath --owner murmur:murmur --directory ${pidfile%/*}
}

# uses the openrc templates for start()/stop()


