# testing/foolsm/foolsm.initd
#!/sbin/openrc-run

command=/usr/sbin/foolsm

if yesno "$use_supervisor"; then
    command_args_foreground="-f"
    supervisor=supervise-daemon
else
    pidfile="/run/$RC_SVCNAME.pid"
    command_args="-p $pidfile"
fi
