# testing/betula/betula.initd
#!/sbin/openrc-run

name=betula
description="Self-hosted personal link collection manager"

command=/usr/bin/betula
command_user=betula:betula
command_args="${datadir}/links ${opts}"

depend() {
	use logger dns
	need net
	after firewall
}

start_pre() {
	checkpath -d -o $command_user "$datadir"
}
