# community/grommunio-common/valkey@grommunio.initd
#!/sbin/openrc-run

# dummy init file which is overwritten with symlink to valkey
