# community/bluez-alsa/bluealsa.initd
#!/sbin/openrc-run
depend() {
	need bluetooth
}
command="/usr/bin/$RC_SVCNAME"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
