# main/open-vm-tools-grsec/vmware-modules.initd
#!/sbin/runscript

MODULES="vmblock vmci vmhgfs vmsync vmxnet vsock vmw_balloon"

depend() {
        before open-vm-tools
}


start() {
        ebegin "Loading VMware kernel modules"
        for mod in $MODULES; do
		modprobe -q $mod
	done
	eend $?
}

stop() {
        ebegin "Unloading VMware kernel modules"
        for mod in $MODULES; do
		rmmod $mod
	done
        eend $?
}

