# testing/netatalk/netatalk.initd
#!/sbin/openrc-run

extra_started_commands="reload"
description_reload="Reload configuration file"

command="/usr/sbin/netatalk"

depend() {
	need net
	use logger dns
	after firewall
}

reload() {
	ebegin "Reloading $name"

	start-stop-daemon --signal SIGHUP --name $command
	eend $?
}
