# testing/netatalk/netatalk.initd
#!/sbin/openrc-run

name=netatalk
command="/usr/sbin/$name"
command_args="-F /etc/afp.conf"
command_background="yes"

pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}

stop() {
	ebegin "Stopping $name"
	start-stop-daemon -K -p /var/lock/netatalk
	eend $?
}
