# testing/empede/empede.initd
#!/sbin/openrc-run

command="/usr/bin/empede"
command_user="empede:empede"
command_background=true
directory="/usr/lib/empede"
pidfile="/run/$RC_SVCNAME.pid"

output_log="/var/log/empede.log"

depend() {
	need net
}

start_pre() {
	checkpath -f -m 0644 -o "$command_user" "$output_log"
}
