# community/irtt/irtt.initd
#!/sbin/openrc-run

description="Isochronous Round-Trip Tester"
name=irtt
command=/usr/bin/irtt
command_user=irtt:irtt
command_background=yes
pidfile=/run/irtt/irtt.pid
command_args="server --syslog=local:"

depend() {
	need net
	use logger
	after firewall
}

start_pre() {
	checkpath --directory --mode 0775 $(dirname $pidfile)
}
