# community/irtt/irtt.initd
#!/sbin/openrc-run

: ${USER:="irtt"}
: ${GROUP:="irtt"}

name="${RC_SVCNAME}"
description="Isochronous Round-Trip Tester"
command=/usr/bin/irtt
command_args="server --syslog=local: ${IRTT_ARGS}"
command_background=yes
command_user="${USER}:${GROUP}"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	use logger
	after firewall
}
