# community/ipfixcol2/ipfixcol2.initd
#!/sbin/openrc-run

name="IPFIXcol2"
description="High-performance NetFlow v5/v9 and IPFIX collector"

: ${cfgfile:="/etc/ipfixcol2/startup.xml"}
: ${command_user:="ipfixcol:ipfixcol"}

command="/usr/bin/ipfixcol2"
command_args="-c $cfgfile $command_args"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

required_files="$cfgfile"

depend() {
	need net
	after firewall
}
