# community/tiny-dfr/tiny-dfr.initd
#!/sbin/openrc-run

pidfile=/run/${RC_SVCNAME}.pid
command=/usr/bin/tiny-dfr
command_background=true

depend() {
	need localmount
	after modules-load
}


