# testing/deviced/deviced.initd
#!/sbin/openrc-run

description="deviced"

command="/usr/bin/deviced"

supervisor=supervise-daemon

depend() {
	need dbus
}
