# community/postfwd/postfwd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=postfwd
description="Postfix firewall daemon"

command=/usr/sbin/postfwd
command_args="--file ${POSTFWD_CONFIG}
	--interface=${POSTFWD_LISTEN}
	--port=${POSTFWD_PORT}
	--user=${POSTFWD_USER}
	--group=${POSTFWD_GROUP}
	${POSTFWD_OPTS}
	"
command_args_foreground="--nodaemon"

depend() {
	need net
}
