# testing/prometheus-rethinkdb-exporter/prometheus-rethinkdb-exporter.initd
#!/sbin/openrc-run

description="Prometheus RethinkDB Exporter"
supervisor=supervise-daemon
command=/usr/bin/prometheus-rethinkdb-exporter

depend() {
    need net
}
