# community/urfkill/urfkill.initd
#!/sbin/openrc-run

description="Handle rfkill events in userspace"
command="/usr/libexec/urfkilld"

depend() {
	after ofono
}
