# community/accel-ppp/accel-pppd.initd
#!/sbin/openrc-run

pidfile="/var/run/accel-pppd.pid"
command="/usr/sbin/$SVCNAME"
command_args="${ACCEL_PPPD_OPTS:--c /etc/accel-ppp.conf -p $pidfile}"
command_background=1
extra_started_commands="reload"

depend() {
        need net
        use dns logger
}

reload() {
  ebegin "Reloading ${RC_SVCNAME}"
  start-stop-daemon --signal USR1 --pidfile "${pidfile}"
  eend $?
}
