# testing/wpa_actiond/wpa_actiond.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name=wpa_actiond
description="Daemon that connects to WPA Supplicant an dhandles events"

command=/usr/sbin/wpa_actiond
command_args="-i $WPA_ACTIOND_INTERFACE -a $WPA_ACTIOND_SCRIPT $WPA_ACTIOND_OPTS"
command_args_foreground="-n"

depend(){
	need net wpa_supplicant
}
