# main/dahdi-tools/dahdi-tools.initd
#!/sbin/runscript

conf=/etc/dahdi/system.conf

depend() {
	before asterisk
	after hwdrivers modules
	keyword novserver
}

start() {
	ebegin "Starting dahdi"
	/usr/sbin/dahdi_cfg
	eend $?
}

stop() {
	ebegin "Stopping dahdi"
	/usr/sbin/dahdi_cfg -s
	eend $?
}
