# community/at/atd.initd
#!/sbin/openrc-run

description="A command scheduling daemon"

command="/usr/sbin/atd"
command_args="
	-f
	${load_average:+"-l $load_average"}
	${batch_interval:+"-b $batch_interval"}
	$command_args
	"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need clock
	use logger
	after localmount
}
