# testing/kea/kea-ctrl-agent.initd
#!/sbin/openrc-run
  
name="kea-ctrl-agent"
command="/usr/sbin/kea-ctrl-agent"
command_args="-c /etc/kea/kea-ctrl-agent.conf"
command_background="yes"
pidfile="/run/kea-ctrl-agent.pid"

depend() {
	need net
	after firewall
}
