# testing/kea/kea-dhcp-ddns.initd
#!/sbin/openrc-run
  
name="kea-dhcp-ddns"
command="/usr/sbin/kea-dhcp-ddns"
command_args="-c /etc/kea/kea-dhcp-ddns.conf"
command_background="yes"
pidfile="/run/kea-dhcp-ddns.pid"

depend() {
	need net
	after firewall
}
