# testing/kea/kea-dhcp4.initd
#!/sbin/openrc-run
  
name="kea-dhcp4"
command="/usr/sbin/kea-dhcp4"
command_args="-c /etc/kea/kea-dhcp4.conf"
command_background="yes"
pidfile="/run/kea-dhcp4.pid"

depend() {
	need net
	after firewall
}
