# testing/kea/kea-dhcp6.initd
#!/sbin/openrc-run
  
name="kea-dhcp6"
command="/usr/sbin/kea-dhcp6"
command_args="-c /etc/kea/kea-dhcp6.conf"
command_background="yes"
pidfile="/run/kea-dhcp6.pid"

depend() {
	need net
	after firewall
}
