# testing/kea/kea.initd.in
#!/sbin/openrc-run

name="@@NAME@@"
command="/usr/sbin/@@NAME@@"
command_args="-c /etc/kea/@@NAME@@.conf"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
