# community/mako/mako.initd
#!/sbin/openrc-run

command=/usr/bin/mako
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
	need dbus
}
