# main/ndisc6/rdnssd.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/sbin/rdnssd
command_args="-H /etc/rdnssd/resolvconf"
command_args_foreground="-f"

start_pre() {
        checkpath -d -m 0755 -o nobody:nobody -q /run/rdnssd
}
