# community/dhcp-probe/dhcp-probe.initd
#!/sbin/openrc-run

depend() {
  need net
  use syslog

}

start() {
	ebegin "Starting DHCP Probe"
	start-stop-daemon --start --quiet --background --make-pidfile \
	--pidfile /var/run/dhcp_probe.pid --exec /usr/sbin/dhcp_probe \
	-- $PARAMS
	eend $?
}

stop() {
	ebegin "Stopping DHCP Probe"
	start-stop-daemon --stop --quiet --pidfile /var/run/dhcp_probe.pid
	eend $?
}
