# community/lxdm/lxdm.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Display Manager"
description="Lightweight Display Manager from LXDM"

command=/usr/sbin/lxdm

depend() {
	provide display-manager
	need localmount
}
