# testing/birdwatcher/birdwatcher.initd
#!/sbin/openrc-run

description="A JSON API for BIRD"
name=birdwatcher
command=/usr/bin/birdwatcher
command_background=yes
pidfile=/run/birdwatcher/birdwatcher.pid

depend() {
	need bird
	use logger
	after firewall
}

start_pre() {
	checkpath --directory --mode 0775 $(dirname $pidfile)
}
