# testing/npd6/npd6.initd
#!/sbin/openrc-run

command="/usr/bin/npd6"
command_args="--foreground $ARGS"
command_background=true
pidfile="/run/npd6.pid"

depend() {
	need net
}
