# testing/htpdate/htpdate.initd
#!/sbin/openrc-run

pidfile=/run/htpdate.pid
command=/usr/sbin/htpdate
command_args="-D -i $pidfile $HTP_OPTIONS $HTP_SERVERS"

depend() {
        need net
}
