# main/hvtools/hv_fcopy_daemon.initd
#!/sbin/openrc-run

name="Hyper-V daemon: hv_fcopy_daemon"
command=/usr/bin/hv_fcopy_daemon

depend() {
	use clock logger
	need localmount
}

