# main/hvtools/hv_vss_daemon.initd
#!/sbin/openrc-run

name="Hyper-V daemon: hv_vss_daemon"
command=/usr/bin/hv_vss_daemon

depend() {
	use clock logger
	need dev
}

