# community/kanshi/kanshi.initd
#!/sbin/openrc-run

description="Dynamic output configuration deamon for Wayland"

command="/usr/bin/kanshi"
supervisor=supervise-daemon

extra_started_commands="reload"

reload() {
	ebegin "Reloading $name configuration"
	/usr/bin/kanshictl reload
	eend $?
}
