# testing/espeakup/espeakup.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="espeakup"
command=/usr/bin/espeakup
command_args="$ESPEAKUP_OPTS"
command_args_foreground="-d"
