# main/bacula/bacula-fd.initd
#!/sbin/openrc-run

: ${port:=9102}
: ${cfgfile:=/etc/bacula/bacula-fd.conf}
: ${command_user:=root}
: ${command_group:=root}

name="Bacula File Daemon"
command="/usr/sbin/${RC_SVCNAME}"
command_args="-u $command_user -g $command_group -c $cfgfile"
pidfile="/run/bacula/${RC_SVCNAME}.$port.pid"

depend() {
	need net
	use dns
	after firewall
}

