# main/bacula/bacula-dir-init
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/www/viewcvs.gentoo.org/raw_cvs/gentoo-x86/app-backup/bacula/files/3.0.3/bacula-dir-init,v 1.1 2009/11/08 17:14:46 wschlich Exp $

depend() {
	need net
	after firewall
	use dns bacula-fd bacula-sd postgresql mariadb
}

start_pre() {
	checkpath --directory --owner bacula:bacula /var/run/bacula
}

start() {
	ebegin "Starting bacula director"
	start-stop-daemon --start --quiet --exec /usr/sbin/bacula-dir \
		-- ${DIR_OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping bacula director"
	start-stop-daemon --stop --quiet --pidfile /var/run/bacula/bacula-dir.*.pid
	eend $?
}
