# community/netbird/netbird-signal.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/netbird-signal
command_user=netbird-signal:netbird

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -m0775 -o root:netbird /var/log/netbird
	checkpath -d -m0755 -o "$command_user" /var/lib/netbird/signal
}
