# community/netbird/netbird-mgmt.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/netbird-mgmt
command_user=netbird-mgmt:netbird

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -m0775 -o root:netbird /var/log/netbird
	checkpath -d -m0755 -o "$command_user" /var/lib/netbird/mgmt 
}
