# community/confd/confd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

command="/usr/bin/confd"
command_background="yes"

pidfile="/var/run/${SVCNAME}.pid"

depend() {
	need net
	after firewall
}
