# testing/tlsrouter/tlsrouter.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name=tlsrouter
command="/usr/bin/tlsrouter"
command_args="$tlsrouter_opts"
command_user="tlsrouter"
start_stop_daemon_args="--quiet"

depend() {
	need net
	after firewall
}
