# non-free/unifi/unifi_db_prune.cron
#!/bin/sh
mongo --port 27117 < /etc/unifi/mongo_prune.js
