# testing/turnstile/turnstiled.initd
#!/sbin/openrc-run

name=turnstiled
command="/usr/bin/turnstiled"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need seatd
}
