# testing/turnstile/turnstiled.initd
#!/sbin/openrc-run

name=turnstiled
command="/usr/bin/turnstiled"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

if [ "$supervisor" = "supervise-daemon" ]; then
	: ${output_log:=/var/log/turnstiled.log}
	: ${error_log:=/var/log/turnstiled.log}
fi

depend() {
	need seatd
}
