# community/vvmd/vvmd.initd
#!/sbin/openrc-run

description="vvmd daemon"

command="/usr/bin/vvmd"
supervisor=supervise-daemon
