# main/dhcpcd/dhcpcd.initd
#!/sbin/openrc-run

description="DHCP Client Daemon"

command="/sbin/dhcpcd"
command_args="-q ${command_args:-}"
command_args_foreground="-B"
pidfile="/run/dhcpcd.pid"

depend() {
	provide net
	need localmount
	use logger network
	after bootmisc modules
	before dns
}
