# main/dhcpcd/dhcpcd.initd
#!/sbin/openrc-run

name="DHCP Client Daemon"
command="/sbin/dhcpcd"
pidfile="/run/dhcpcd.pid"
command_args="-q ${command_args:-}"

depend() {
	provide net
	need localmount
	use logger network
	after bootmisc modules
	before dns
}
