# community/lua-cqueues-pushy/cqpushy.initd
#!/sbin/openrc-run

depend() {
	after net firewall chrony dbus
}

description="lua-cqueues-pushy"
command="/etc/cqpushy/cqpushy.lua"
command_background="yes"
pidfile="/var/run/cqpushy.pid"
