# community/swaybg/swaybg.initd
#!/sbin/openrc-run

command="/usr/bin/swaybg"
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
