# community/nut/powerfail.initd
#!/sbin/openrc-run
# Copyright 2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description='Signal the UPS to kill power in a power failure condition'

depend() {
	need mount-ro
}

start() {
	if [ "$RC_RUNLEVEL" = "shutdown" -a -f /etc/killpower ] ; then
		ebegin 'Signaling UPS to kill power'
		/usr/sbin/upsdrvctl shutdown
		eend $?
	fi
}
