# community/rtkit/rtkit.initd
#!/sbin/openrc-run

command="/usr/libexec/rtkit-daemon"
command_background="yes"
capabilities="cap_sys_nice,cap_dac_read_search,cap_sys_chroot,cap_setgid,cap_setuid"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need dbus
}
