# testing/icecast/icecast.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting Icecast 2"
	start-stop-daemon --start --exec /usr/bin/icecast \
		-- \
		-b -c ${icecast_conf:-/etc/icecast.xml} > /dev/null
	eend $?
}

stop() {
	ebegin "Stopping Icecast 2"
	start-stop-daemon --stop --exec /usr/bin/icecast
	eend $?
}
