# testing/icecast/icecast.initd
#!/sbin/runscript

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting Icecast 2"
	start-stop-daemon --background --start --make-pidfile --pidfile /var/run/icecast.pid --exec /usr/bin/icecast -- -c /etc/icecast.xml
	eend $?
}

stop() {
	ebegin "Stopping Icecast 2"
	start-stop-daemon --stop --pidfile /var/run/icecast.pid --name icecast
	eend $?
}
