# community/bird-lg-go/bird-lg-go.initd
#!/sbin/openrc-run

name="${RC_SVCNAME}"
pidfile="/run/$name.pid"

command="/usr/bin/bird-lg-go"
: ${command_user:="bird-lg-go:bird-lg-go"}
command_background=true

: ${output_logger="logger -t bird-lg-go -p daemon.info >/dev/null 2>&1"}
: ${error_logger="logger -t bird-lg-go -p daemon.warning >/dev/null 2>&1"}

depend() {
	need net
	use bird bird-lg-proxy
	after firewall
}
