# testing/dante/sockd.initd
#!/sbin/openrc-run
#
# startup script for sockd daemon
#

depend() {
    need net
    after firewall
    provide sockd
    use dns
}

command=/usr/sbin/sockd
pidfile=/var/run/${SVCNAME}.pid
command_args="${SOCKD_OPTS}"

start_pre() {
    $command -V
}
