# community/bonsai/bonsaid.initd
#!/sbin/openrc-run

description="bonsai daemon"

command="/usr/bin/bonsaid"
supervisor=supervise-daemon
