# main/arpon/arpon.initd
#!/sbin/openrc-run

command=/usr/sbin/arpon
command_args="-d $ARPON_OPTS"
pidfile=/var/run/arpon.pid

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f "$pidfile"
}
