# main/arpon/arpon.initd
#!/sbin/runscript

NAME=arpon
DAEMON=/usr/sbin/$NAME

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting ${NAME}"
		start-stop-daemon --start --quiet \
			--pidfile /var/run/${NAME}.pid \
			--exec ${DAEMON} -- -d ${ARPON_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"
		start-stop-daemon --stop --quiet \
			--pidfile /var/run/${NAME}.pid
	eend $?
}
