# community/rest-server/rest-server.initd
#!/sbin/openrc-run

name="rest-server"
command="/usr/bin/$name"
command_args="$REST_SERVER_OPTS --path '${REST_SERVER_PATH}'"
command_background="yes"
command_user="${REST_USER:-restic}"
pidfile="/var/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath --directory "${REST_SERVER_PATH}" \
		--owner "${REST_USER:-restic}:${REST_GROUP:-restic}"
}
