# main/uwsgi/uwsgi.initd
#!/sbin/openrc-run

pidfile="/run/uwsgi/uwsgi.pid"
command="/usr/sbin/uwsgi"
command_args="--daemonize=/var/log/uwsgi.log --emperor /etc/uwsgi/conf.d --emperor-pidfile=$pidfile"

depend() {
        need net
}

start_pre() {
        checkpath --directory --owner uwsgi:uwsgi \
                --mode 0775 /run/uwsgi
}
