# testing/bird/bird.initd
#!/sbin/openrc-run

extra_started_commands="reload"

name=bird
command=/usr/sbin/$name

sock="/var/run/$name.ctl"
command_args="-c /etc/${name}.conf -s $sock"

depend() {
	need net
	use logger
	after firewall
}

start_pre() {
	[ -f "/etc/$name.conf" ] && return 0
	eerror "Please create /etc/$name.conf"
	return 1
}

reload() {
	ebegin "Reloading $name"
	killall -HUP $command
	eend $?
}
