# community/cs-firewall-bouncer/cs-firewall-bouncer.initd
#!/sbin/openrc-run

command=/usr/bin/crowdsec-firewall-bouncer
command_args="-c /etc/crowdsec/bouncers/crowdsec-firewall-bouncer.yaml"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
	after firewall
}
