# testing/ushare/ushare.initd
#!/sbin/runscript

depend() {
  need net
  use syslog
}

start() {
  ebegin "Starting uShare UPnP A/V Media Server"
  start-stop-daemon --start --quiet --background --make-pidfile \
  --pidfile /var/run/ushare.pid --exec /bin/ushare \
  -- $PARAMS
  eend $?
}

stop() {
  ebegin "Stopping uShare UPnP A/V Media Server"
  start-stop-daemon --stop -quiet --pidfile /var/run/ushare.pid
  eend $?
}
