# community/xdg-desktop-portal/xdg-document-portal.initd
#!/sbin/openrc-run

command="/usr/libexec/xdg-document-portal"
supervisor=supervise-daemon
command_args="--replace $command_args"
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
	need dbus
}
