# testing/mtg/mtg.initd
#!/sbin/openrc-run

command="/usr/sbin/mtg"
command_args="run $MTG_OPTS"
command_background=yes
command_user=mtg:mtg
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
