# testing/neard/neard.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="NFC"
command="/usr/libexec/nfc/neard"
command_args_foreground="--nodaemon"

depend() {
	after coldplug
	need dbus
}
