# testing/qbittorrent-nox/qbittorrent-nox.initd
#!/sbin/openrc-run

# init.d file for qbittorrent-nox.

name=qbittorrent-nox
user=qbittorrent
group=qbittorrent
command="/usr/bin/$name"
command_args="-d $qbittorrent_opts"

start_stop_daemon_args="--user $user:$group"

depend() {
	need net
	after firewall
}
