# testing/fwknop/fwknopd.initd
#!/sbin/openrc-run

pidfile="/var/run/fwknop/fwknopd.pid"
command=/usr/sbin/fwknopd
command_args="--pid-file $pidfile $fwknopd_args"

depend() {
	use logger
        after net
}

start_pre() {
	checkpath --directory ${pidfile%/*}
}
