# testing/seaweedfs/seaweedfs.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=seaweedfs
command="/usr/bin/weed"
command_args="${WEED_EXEC} ${WEED_OPTS}"

pidfile="/run/$name.pid"
respawn_delay=5

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f -m 0644 -o root:root /var/lib/seaweedfs/metadata
}

