# testing/thermald/thermald.initd
#!/sbin/openrc-run

command=/usr/sbin/thermald
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
command_args="--no-daemon --dbus-enable"

depend() {
	use clock logger
	need localmount
}
