# testing/thermald/thermald.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="Thermal Monitor"
description="Read temperature sensors and decide best action based on data and user preferences"

command=/usr/sbin/thermald
command_args="--dbus-enable"
command_args_foreground="--no-daemon"

depend() {
	use clock logger
	need localmount
}
