# testing/prometheus-postgres-exporter/postgres-exporter.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

command="/usr/bin/postgres_exporter"
command_args="$ARGS"
command_background="yes"
command_user="prometheus:prometheus"

logdir="/var/log/prometheus"
error_log="$logdir/${SVCNAME}.log"
pidfile="/var/run/${SVCNAME}.pid"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -d -o $user:$group -m755 $logdir
	checkpath -f -o $user:$group -m644 $error_log
}
