# testing/py3-pyatem/openswitcher-proxy.initd
#!/sbin/openrc-run

name="OpenSwitcher Proxy"
description="Reverse proxy for the ATEM protocol"

: ${cfgfile:="/etc/openswitcher/proxy.toml"}

command="/usr/bin/openswitcher-proxy"
command_args="--config $cfgfile"
command_user="openswitcher:openswitcher"

required_files="$cfgfile"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f --owner "$command_user" --mode 0644 $logfile
}
