# testing/speedtest-go/speedtest-go.initd
#!/sbin/openrc-run

name="speedtest-go"
description="HTML5 Speedtest, written in go"

command="/usr/bin/speedtest"
command_args="-c /etc/speedtest-go/settings.toml"
command_user="nobody:nobody"
command_background="true"
pidfile="/var/run/${name}.pid"

depend() {
	need net
}
