# testing/recoll/recollindex.initd
#!/sbin/openrc-run

description="Recoll Index Daemon"

command="/usr/bin/recollindex"
command_background=yes
command_args_foreground="-D -O -m ${recollindex_args}"

error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

