# testing/rasdaemon/rasdaemon.initd
#!/sbin/openrc-run

description="daemon to receive RAS error tracings"

pidfile=/run/"$RC_SVCNAME".pid
command=/usr/sbin/rasdaemon
command_args="--foreground --record"
command_background=true

depend() {
	need localmount sysfs
}

start_post() {
	"${command}" --enable
}

stop_post() {
	"${command}" --disable
}
