# community/pimsync/pimsync.initd
#!/sbin/openrc-run

name="Pimsync"
description="Synchronise calendars and contacts"

supervisor=supervise-daemon
command="/usr/bin/pimsync"
command_args="-v info daemon -r 3"
notify=fd:3
