# testing/gitea-runner/gitea-runner.initd
#!/sbin/openrc-run

description="Gitea CI Runner"
name="Gitea Runner"

: ${cfgfile:="/etc/gitea-runner/config.yaml"}
: ${datadir:="/var/lib/gitea-runner"}
: ${command_user:="gitea-runner"}

command="/usr/bin/gitea-runner"
command_args="daemon --config $cfgfile"
command_background="yes"
directory="$datadir"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	use dns logger
}

start_pre() {
	checkpath -d -o "$command_user" /etc/gitea-runner
	checkpath -d -o "$command_user" "$datadir"

	if ! [ -e "$cfgfile" ]; then
		eerror "Config file $cfgfile doesn't exist."
		eerror "You can generate it with: gitea-runner generate-config,"
		eerror "or use the auto-generated one in /etc/gitea-runner/config.example.yaml"
		return 1
	fi

	if [ "$error_log" ]; then
		output_log="$error_log"
	else
		output_logger="logger -t '${RC_SVCNAME}' -p daemon.info"
		error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
	fi
}

