# testing/conntracct/conntracct.initd
#!/sbin/openrc-run

name=conntracct
description="Conntracct is a tool for extracting network flow information from Linux hosts"
command=/usr/bin/conntracct
command_args="$ARGS"
command_user="conntracct:conntracct"
supervisor="supervise-daemon"
capabilities="^cap_sys_admin,^cap_net_admin,^cap_dac_override,^cap_sys_resource"

depend() {
        need net
}
