# community/coturn/turnserver.initd
#!/sbin/openrc-run

description="TURN and STUN Server"
name=turnserver
command=/usr/bin/turnserver
command_user=turnserver
pidfile=/run/turnserver/turnserver.pid
command_args="--no-stdout-log --daemon --pidfile $pidfile"

start_pre() {
	checkpath --directory --owner turnserver:turnserver --mode 0775 \
		/run/turnserver /var/log/turnserver
}
