# testing/geodns/geodns.initd
#!/sbin/openrc-run

command=/usr/bin/geodns
command_args="-privatedebug -config /etc/geodns/ -logfile /var/log/geodns \
 -port ${PORT-53} -http :${HTTP_PORT-8053}"
supervisor=supervise-daemon

