# testing/seatd/seatd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name=$RC_SVCNAME
description="Seat management daemon"

command="seatd"
command_args="-g video"

supervise_daemon_args="--stderr /var/log/seatd.log"
