# community/caddy/caddy.initd
#!/sbin/openrc-run

name="Caddy webserver"
command="/usr/sbin/caddy"
command_args="$CADDY_OPTS"
pidfile="/var/run/caddy.pid"
command_background=yes
start_stop_daemon_args="--user caddy --group caddy"

depend() {
	need net localmount
	after firewall
}
