# main/haveged/haveged.initd
#!/sbin/openrc-run

description="Entropy harvesting daemon using CPU timings"

command="/usr/sbin/haveged"
# $HAVEGED_OPTS is here for backward compatiblity only
command_args="--Foreground ${command_args:-$HAVEGED_OPTS}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need localmount
	after urandom
	provide entropy
}
