# main/haveged/haveged.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="RNG Daemon"
description="haveged Random Number Generator"

command=/usr/sbin/haveged
command_args="$HAVEGED_OPTS"
command_args_foreground="--Foreground"

depend() {
	need localmount
	after urandom
	provide entropy
}
