# main/haveged/haveged.initd
#!/sbin/openrc-run

command="/usr/sbin/haveged"
command_args="$HAVEGED_OPTS"
command_background="yes"

pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need localmount
	after urandom
	provide entropy
}
