# community/addrwatch/addrwatch-syslog.initd
#!/sbin/openrc-run

name="addrwatch-syslog"
description="addrwatch syslog output module"

: ${command_user:="addrwatch"}
: ${error_logger="logger -t addrwatch-syslog -p daemon.err >/dev/null 2>&1"}

command="/usr/bin/addrwatch_syslog"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need addrwatch
}
