# community/iipsrv/iipsrv.initd
#!/sbin/openrc-run

description="IIPImage server"
supervisor=supervise-daemon

command="/usr/bin/iipsrv"
command_args="--bind $HOST:$PORT"
command_user="iipsrv:wheel"

start_pre() {
	# Make sure logfile is writable
	touch $LOGFILE
	chown iipsrv:wheel $LOGFILE
}
