# community/dsme/dsme.initd
#!/sbin/openrc-run

name="DSME"
description="Device State Management Entity"

command=/usr/sbin/dsme
command_args="-- -p /usr/lib/dsme/startup.so"

depend() {
	need dbus
}
