# testing/chasquid/chasquid.initd
#!/sbin/openrc-run

name="shasquid"
description="chasquid mail daemon (service)"

: ${command_user:="chasquid"}

command="/usr/bin/chasquid"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
