# main/macifrename/macifrename.initd
#!/sbin/openrc-run

depend() {
        before net
        keyword -vserver -lxc
}

start() {
	ebegin "Renaming network devices"

	/usr/sbin/macifrename ${MACIFRENAME_OPTS}

	eend $?
}
