# community/xrdp/xrdp.initd
#!/sbin/openrc-run

command=/usr/sbin/$SVCNAME
pidfile=/var/run/$SVCNAME.pid

depend() {
	need net
	after firewall
}
