# main/mumble/murmur.initd
#!/sbin/openrc-run

depends() {
	needs localmount
}

description="Mumble VoIP Server"
command="/usr/bin/murmurd"
command_args="-ini /etc/murmur.ini"
pidfile="/var/run/murmur/murmur.pid"

start_pre() {
	checkpath --owner murmur:murmur --directory ${pidfile%/*}
}

# uses the openrc templates for start()/stop()


