# community/tigervnc/vncserver.initd
#!/sbin/openrc-run

supervisor=supervise-daemon

name="TigerVNC server"

command="/usr/bin/vncserver"
command_args="$VNC_DISPLAY"
command_user="$VNC_USER:$VNC_GROUP"

output_log="/var/log/vncserver.log"
error_log="/var/log/vncserver.log"

depend() {
	need net
	after firewall
}

start_pre() {
	checkpath -f -m 644 -o "$command_user" "$output_log"
}
