# community/elogind/elogind.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="System login manager"
description="System service that manages user logins"

command=/usr/libexec/elogind/elogind

depend() {
	need dbus

	# Make sure we start before any other display manager
	before xdm lxdm sddm gdm
}
