# testing/matterbridge/matterbridge.initd
#!/sbin/openrc-run

name=matterbridge
description="Bridges between a growing number of chat protocols"

command="/usr/bin/matterbridge"
command_args="$MATTERBRIDGE_OPTS"
command_background=yes
command_user=matterbridge:matterbridge
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need net
	after firewall
}
