# testing/upterm/uptermd.initd
#!/sbin/openrc-run

name=uptermd
description="Secure terminal sharing deamon"
command="/usr/bin/uptermd"
command_args="${command_args:---ssh-addr 0.0.0.0:2222}"
command_background="yes"

depend() {
	need net
	after firewall
}
