# testing/upterm/uptermd.initd
#!/sbin/openrc-run

: ${command_user:="uptermd:uptermd"}

name=uptermd
description="Secure terminal sharing deamon"
command="/usr/bin/uptermd"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
capabilities="^cap_net_bind_service"

depend() {
	need net
	after firewall
}
