# community/wayland-pipewire-idle-inhibit/wayland-pipewire-idle-inhibit.initd
#!/sbin/openrc-run

name="wayland-pipewire-idle-inhibit"
description="wayland-pipewire-idle-inhibit"

supervisor=supervise-daemon
command=/usr/bin/wayland-pipewire-idle-inhibit
output_logger="logger -t '$RC_SVCNAME' -p daemon.info"

. /etc/user/conf.d/dbus

depend() {
	need dbus pipewire
}
