# community/nbd/nbd-server.initd
#!/sbin/openrc-run

command=/usr/bin/nbd-server
pidfile=/run/nbd-server.pid
command_args="-p ${pidfile}"
