# testing/spotifyd/spotifyd.initd
#!/sbin/openrc-run

command=/usr/bin/spotifyd
command_args="${spotifyd_args}"
pidfile="/run/spotifyd.pid"
command_background=true
start_stop_daemon_args="--user spotifyd --group audio"

depends() {
	need localmount net
	before firewall
}
