# community/opendkim/opendkim.initd
#!/sbin/openrc-run

owner=opendkim
pidfile=/run/opendkim/opendkim.pid
cfgfile=/etc/opendkim/opendkim.conf
command=/usr/sbin/opendkim
command_args="$command_args -u $owner -P $pidfile"
required_files="$cfgfile"

depend() {
	need net
	before mta
}

start_pre() {
	checkpath -d -o $owner ${pidfile%/*}
}
