# testing/zabbix-agent2-alpine/zabbix-agent2.initd
#!/sbin/openrc-run

command=/usr/sbin/zabbix_agent2
supervisor="supervise-daemon"
command_user="zabbix:zabbix"

start_pre() {
    checkpath --directory --owner zabbix:zabbix /run/zabbix/
    checkpath --file --owner zabbix:zabbix /var/log/zabbix/zabbix_agent2.log
    checkpath --directory --owner zabbix:zabbix /var/log/zabbix
}
