# community/mimir/mimir.initd
#!/sbin/openrc-run

configfile="${MIMIR_CONF:-/etc/mimir/mimir.yaml}"

command=/usr/bin/mimir
command_args="-config.file=${configfile}
	${MIMIR_OPTS}"
command_background=yes
command_user=mimir:grafana
make_pidfile=yes
pidfile=/var/run/mimir.pid
