# testing/tuwunel/tuwunel.initd
#!/sbin/openrc-run
supervisor="supervise-daemon"

name=$RC_SVCNAME
description="Featureful Matrix homeserver"

command="/usr/bin/tuwunel"
: ${command_args:="--config $cfgfile"}
command_user="tuwunel:tuwunel"

depend() {
        need net
        after firewall
}

start_pre() {
        checkpath -d -o $command_user $rocksdb_dir
}
